# bastion-api.service
[Unit]
Description=Bastion Headless API Server (CCC)
After=network.target

[Service]
Type=simple
User=ccc
WorkingDirectory=/home/ccc/ccc
EnvironmentFile=-/home/ccc/ccc/.env
Environment=PYTHONPATH=/home/ccc/ccc
ExecStart=/home/ccc/ccc/.venv/bin/uvicorn apps.bastion.api:app --host 0.0.0.0 --port 8003
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
