# contents/standalone/lab/secuops/week07.yaml
lab_id: standalone-secuops-week07
title: osquery 호스트 가시화 (4 호스트) + R/B/P 헌팅
version: standalone-non-ai
course: secuops
week: 7
description: |
  bastion / fw / ips / web 4 호스트의 osquery 5.x 로 OS-as-SQL 호스트 가시화 실습.
  processes / users / listening_ports / suid_bin / file_events 핵심 테이블 + 헌팅
  쿼리 + Red/Blue/Purple (Red 가 의심 프로세스 spawn → Blue 가 헌팅 → Purple 가
  FIM 권장).
difficulty: medium
duration_minutes: 150
prerequisites: [W04-W06]
objectives:
  - "osqueryi 대화형 + JSON 출력"
  - "5 핵심 테이블 (processes/users/listening_ports/suid_bin/hash) 쿼리"
  - "헌팅 쿼리 (on_disk=0, world-writable, cron backdoor)"
  - "FIM 설정 (file_events) + 변경 감지"
  - "R/B/P — 의심 프로세스 → 헌팅 → FIM 권장"
pass_threshold: 0.7
steps:

  - order: 1
    instruction: |
      ## 🎯 MISSION

      4 호스트 (bastion / fw / ips / web) 의 **osqueryi 가용성** + **kernel/os 정보**
      가시화.

      ## 💻 실행 명령

      ```bash
      for h in bastion fw ips web; do
        echo "=== 6v6-$h ==="
        ssh 6v6-$h 'sudo osqueryi --json "SELECT version FROM osquery_info;"'
      done
      ```

      ## 📋 확인 사항

      - 4 호스트 모두 `[{"version":"5.x.x"}]` 형식 출력
      - 5.x 버전 (pkg.osquery.io 공식 배포)

      ## 💡 왜 이걸 하는가?

      W07 의 모든 step 의 전제 — 4 호스트에 osquery 설치 여부. 미설치 시 lecture 의
      "Dockerfile 에 osquery 추가" 절차 재실행 필요.

      ## 🔍 결과 해석

      - **정상**: 4 호스트 모두 version 응답
      - **비정상**: 일부 호스트 `command not found` → 해당 컨테이너 Dockerfile 의 osquery
        패키지 누락 → 재빌드
    hint: ssh 6v6-X 'sudo osqueryi --json "SELECT version FROM osquery_info;"'
    category: survey
    points: 10
    answer: |
      for h in bastion fw ips web; do
        echo "=== 6v6-$h ==="
        ssh 6v6-$h 'sudo osqueryi --json "SELECT version FROM osquery_info;"'
      done
    answer_detail: 4 호스트 osquery 5.x 가용.
    verify:
      type: output_contains
      expect: version
      field: stdout
      semantic:
        intent: "4 호스트 osquery 가용성"
        success_criteria: ["4 호스트 응답", "version 5.x"]
        acceptable_methods:
          - "**[OSS] osqueryi --json (JSON 출력, 자동화 친화)**"
          - "**[OSS] osqueryi --csv (spreadsheet)**"
          - "**[OSS] osqueryi 대화형 (.tables / .schema)**"
          - "**[OSS] osqueryd daemon 모드 (scheduled query)**"
          - "**[OSS] kernel_info / os_version / platform_info**"
          - "**[OSS] systemctl status osqueryd**"
          - "**[OSS] /var/log/osquery/osqueryd.results.log**"
          - "**[OSS] osqueryctl status / start / stop / restart**"
    target_vm: bastion

  - order: 2
    instruction: |
      ## 🎯 MISSION (5 핵심 테이블)

      web 호스트에서 **5 핵심 테이블** (processes / users / listening_ports / suid_bin
      / hash) 쿼리.

      ## 💻 실행 명령

      ```bash
      # 1. 프로세스 (apache2 / sshd)
      ssh 6v6-web 'sudo osqueryi --json "SELECT pid, name, cmdline, uid FROM processes WHERE name IN ('"'"'apache2'"'"','"'"'sshd'"'"');"'

      # 2. 사용자 (uid >= 1000)
      ssh 6v6-web 'sudo osqueryi --json "SELECT username, uid, gid, shell FROM users WHERE uid >= 1000;"'

      # 3. listening port + process 매핑
      ssh 6v6-web 'sudo osqueryi --json "SELECT l.port, p.name FROM listening_ports l JOIN processes p ON l.pid=p.pid;"'

      # 4. SUID binary top 10
      ssh 6v6-web 'sudo osqueryi --json "SELECT path FROM suid_bin LIMIT 10;"'

      # 5. /etc/passwd 해시
      ssh 6v6-web 'sudo osqueryi --json "SELECT path, sha256 FROM hash WHERE path = '"'"'/etc/passwd'"'"';"'
      ```

      ## 📋 확인 사항

      - 명령 1: apache2 + sshd 프로세스 (cmdline 포함)
      - 명령 2: ccc 사용자 (uid 1000)
      - 명령 3: 22 sshd / 80 apache2 / 443 apache2 매핑
      - 명령 4: /usr/bin/sudo / /usr/bin/su 등 SUID binary
      - 명령 5: /etc/passwd 의 sha256 해시 (변경 감지 기준선)

      ## 💡 왜 이걸 하는가?

      osquery 의 강점 — OS 의 모든 정보를 SQL 로 추상화. JOIN / WHERE / LIMIT 같은
      표준 SQL 기능 활용 가능.

      ## 🔍 결과 해석

      - **정상**: 5 쿼리 모두 정상 응답
      - **비정상**: hash 가 권한 부족 → sudo 필요
    hint: osqueryi --json + processes/users/listening_ports/suid_bin/hash
    category: monitoring
    points: 15
    answer: |
      ssh 6v6-web 'sudo osqueryi --json "SELECT pid, name, cmdline, uid FROM processes WHERE name IN ('"'"'apache2'"'"','"'"'sshd'"'"');"'
      ssh 6v6-web 'sudo osqueryi --json "SELECT username, uid, gid, shell FROM users WHERE uid >= 1000;"'
      ssh 6v6-web 'sudo osqueryi --json "SELECT l.port, p.name FROM listening_ports l JOIN processes p ON l.pid=p.pid;"'
      ssh 6v6-web 'sudo osqueryi --json "SELECT path FROM suid_bin LIMIT 10;"'
      ssh 6v6-web 'sudo osqueryi --json "SELECT path, sha256 FROM hash WHERE path = '"'"'/etc/passwd'"'"';"'
    answer_detail: 5 핵심 테이블 쿼리. JOIN 활용 + LIMIT + WHERE 표준 SQL.
    verify:
      type: output_contains
      expect: apache2
      field: stdout
      semantic:
        intent: "5 핵심 테이블 쿼리"
        success_criteria: ["processes 결과", "users 결과", "listening_ports JOIN"]
        acceptable_methods:
          - "**[OSS] SELECT FROM processes** (50+ 컬럼)"
          - "**[OSS] users + uid >= 1000 (system 계정 제외)**"
          - "**[OSS] listening_ports JOIN processes — daemon-port 매핑**"
          - "**[OSS] suid_bin — 권한 상승 잠재 위험**"
          - "**[OSS] hash — 파일 무결성 (sha256)**"
          - "**[OSS] authorized_keys — SSH 인증 키**"
          - "**[OSS] shell_history — bash history**"
          - "**[OSS] sudoers — sudo 권한 audit**"
    target_vm: web

  - order: 3
    instruction: |
      ## 🎯 MISSION (헌팅 쿼리 4종)

      4 호스트에서 **헌팅 쿼리** 실행 — on_disk=0 / world-writable / unusual cron /
      비표준 shell 사용자.

      ## 💻 실행 명령

      ```bash
      # 1. on_disk = 0 process (fileless malware 의심)
      ssh 6v6-web 'sudo osqueryi --json "SELECT pid, name, path FROM processes WHERE on_disk = 0;"'

      # 2. world-writable file in /tmp
      ssh 6v6-web 'sudo osqueryi --json "SELECT path, mode FROM file WHERE path LIKE '"'"'/tmp/%'"'"' AND mode = '"'"'0777'"'"' LIMIT 5;"'

      # 3. cron backdoor 의심
      ssh 6v6-web 'sudo osqueryi --json "SELECT command, path FROM crontab LIMIT 10;"'

      # 4. 비표준 shell 사용자 (bash 외)
      ssh 6v6-web 'sudo osqueryi --json "SELECT username, shell FROM users WHERE shell NOT IN ('"'"'/bin/false'"'"', '"'"'/usr/sbin/nologin'"'"', '"'"'/bin/sh'"'"', '"'"'/bin/bash'"'"');"'
      ```

      ## 📋 확인 사항

      - 명령 1: 빈 결과 `[]` (정상 환경에 fileless malware 없음)
      - 명령 2: 빈 결과 또는 system 임시 파일
      - 명령 3: /etc/crontab 의 표준 cron entry
      - 명령 4: 빈 결과 (모든 사용자가 표준 shell)

      ## 💡 왜 이걸 하는가?

      Threat Hunting 의 표준 쿼리. malware 행위 (fileless / cron backdoor / 비인가
      shell) 검출. 정기 cron 으로 자동 실행.

      ## 🔍 결과 해석

      - **정상**: 4 쿼리 모두 빈 결과 또는 알려진 정상 entry
      - **의심**: on_disk=0 process 발견 → 즉시 incident response 시작 (lsof + memory dump)
    hint: osquery 헌팅 쿼리
    category: hunting
    points: 15
    answer: |
      ssh 6v6-web 'sudo osqueryi --json "SELECT pid, name, path FROM processes WHERE on_disk = 0;"'
      ssh 6v6-web 'sudo osqueryi --json "SELECT path, mode FROM file WHERE path LIKE '"'"'/tmp/%'"'"' AND mode = '"'"'0777'"'"' LIMIT 5;"'
      ssh 6v6-web 'sudo osqueryi --json "SELECT command, path FROM crontab LIMIT 10;"'
      ssh 6v6-web 'sudo osqueryi --json "SELECT username, shell FROM users WHERE shell NOT IN ('"'"'/bin/false'"'"', '"'"'/usr/sbin/nologin'"'"', '"'"'/bin/sh'"'"', '"'"'/bin/bash'"'"');"'
    answer_detail: 4 헌팅 쿼리 정상 환경에서 모두 빈 결과 또는 알려진 entry.
    verify:
      type: output_contains
      expect: '['
      field: stdout
      semantic:
        intent: "헌팅 쿼리 4종 baseline"
        success_criteria: ["4 쿼리 실행"]
        acceptable_methods:
          - "**[OSS] on_disk = 0 (fileless)**"
          - "**[OSS] world-writable file**"
          - "**[OSS] crontab + at + systemd timer**"
          - "**[OSS] authorized_keys + 비표준 shell**"
          - "**[OSS] SUID binary 변경**"
          - "**[OSS] capabilities (getcap)**"
          - "**[OSS] kernel module — kernel_modules**"
          - "**[OSS] startup_items / launch_agents (Linux/Mac)**"
    target_vm: web

  - order: 4
    instruction: |
      ## 🎯 MISSION (R/B/P)

      **Red** 가 attacker 에서 의심 프로세스 (nc listener) 시뮬 → **Blue** 가 osquery
      로 처리 추적 → **Purple** 가 FIM (file_events) 권장.

      ## 💻 실행 명령

      ```bash
      # Red — attacker 에서 nc listener 시뮬 (학습 환경)
      ssh 6v6-attacker 'timeout 5 nc -l -p 9999 &' 2>/dev/null &
      sleep 1

      # Blue — attacker 의 osquery 로 nc 프로세스 추적 (attacker 에도 osquery 설치되어 있다면)
      # 실제로는 web 의 정상 process 가시화로 대체
      ssh 6v6-web 'sudo osqueryi --json "SELECT pid, name, cmdline, start_time FROM processes WHERE name LIKE '"'"'%nc%'"'"' OR name = '"'"'sshd'"'"' LIMIT 5;"'

      # Purple — FIM 권장 (시뮬)
      cat <<'EOF'
      ## FIM 권장 osquery.conf patch
      {
        "file_paths": {
          "sensitive": ["/etc/passwd", "/etc/shadow", "/etc/sudoers"]
        },
        "schedule": {
          "file_changes": {
            "query": "SELECT * FROM file_events WHERE category='sensitive';",
            "interval": 30
          }
        }
      }
      EOF
      ```

      ## 📋 확인 사항

      - Blue: web 에서 sshd 또는 nc 프로세스 가시화
      - Purple: file_events 설정 patch 작성 시뮬 (실 적용은 osqueryd 재시작 필요)

      ## 💡 왜 이걸 하는가?

      Purple Team — Red 가 흘린 의심 프로세스를 Blue 가 osquery 로 추적, Purple 이
      FIM 으로 사전 차단 권장.

      ## 🔍 결과 해석

      - **정상**: Blue 쿼리 결과 + Purple FIM patch
      - **운영**: file_events 의 realtime=yes 설정 → inotify 기반 즉시 감지
    hint: R/B/P + osquery + FIM 권장
    category: redblue
    points: 25
    answer: |
      ssh 6v6-web 'sudo osqueryi --json "SELECT pid, name, cmdline, start_time FROM processes WHERE name = '"'"'sshd'"'"' OR name LIKE '"'"'%apache%'"'"' LIMIT 5;"'
      echo "Purple FIM 권장:"
      cat <<'EOF'
      {
        "file_paths": {
          "sensitive": ["/etc/passwd", "/etc/shadow", "/etc/sudoers"]
        },
        "schedule": {
          "file_changes": {
            "query": "SELECT * FROM file_events WHERE category='sensitive';",
            "interval": 30
          }
        }
      }
      EOF
    answer_detail: R/B/P — 의심 프로세스 추적 + FIM 권장 patch.
    verify:
      type: output_contains
      expect: file_events
      field: stdout
      semantic:
        intent: "R/B/P + FIM 권장"
        success_criteria: ["프로세스 추적", "FIM patch"]
        acceptable_methods:
          - "**[OSS] processes 의 start_time + on_disk 컬럼**"
          - "**[OSS] file_events 테이블 (FIM 결과)**"
          - "**[OSS] osquery.conf 의 file_paths + schedule**"
          - "**[OSS] interval 30s — inotify realtime 모드**"
          - "**[OSS] wildcard % vs %% (재귀)**"
          - "**[OSS] Wazuh agent 의 localfile 로 osqueryd.results.log ship (W10)**"
          - "**[OSS] sysmon-for-linux + osquery 보완 (W11)**"
          - "**[OSS] fleet 솔루션 (kolide/doorman) 대규모 운영**"
    target_vm: web

  - order: 5
    instruction: |
      ## 🎯 MISSION

      W07 종합 보고서 + 헌팅 쿼리 5개 + W09 (Wazuh 통합) 준비도.

      ## 💻 산출물

      ```markdown
      # W07 osquery 호스트 가시화 보고서
      ## 1. 4 호스트 osquery 가용성
      ## 2. 5 핵심 테이블 + 4 헌팅 쿼리 결과
      ## 3. R/B/P (의심 프로세스 + FIM 권장)
      ## 4. W09-W10 (Wazuh + osquery 통합) 준비도
      ## 5. 운영 권장 (osqueryd schedule / fleet / Wazuh 통합)
      ```
    hint: 마크다운
    category: synthesis
    points: 10
    answer: 위 형식
    answer_detail: 종합 보고서
    verify:
      type: output_contains
      expect: osquery
      field: stdout
      semantic:
        intent: "W07 종합"
        success_criteria: ["5 섹션"]
        acceptable_methods:
          - "**[관리적]** 4 host osquery 가용성"
          - "**[관리적]** NIST DE.CM-1 매핑 (continuous monitoring)"
          - "**[관리적]** ISMS-P 2.5 (인증·권한) 매핑"
          - "**[관리적]** MITRE D3-PA (Process Analysis)"
          - "**[관리적]** 운영 권장 — schedule interval + logrotate"
          - "**[관리적]** Wazuh 통합 path (W10)"
          - "**[관리적]** sysmon 보완 (W11)"
          - "**[관리적]** fleet 솔루션 대규모"
    target_vm: web
