# contents/standalone/lab/secuops/week08.yaml
lab_id: standalone-secuops-week08
title: 중간고사 — 5 시나리오 종합 실기 (90분, 100점)
version: standalone-non-ai
course: secuops
week: 8
description: |
  W01-W07 의 종합 평가. 5 시나리오 × 20점. fw nftables + Suricata 룰 + ModSec 공격
  시뮬 + osquery 헌팅 + 통합 R/B/P 분석.
difficulty: hard
duration_minutes: 90
prerequisites: [W01-W07 통과]
objectives:
  - "5 시나리오 모두 본인 손으로 수행"
  - "각 시나리오 결과 (출력 + 분석) 제출"
pass_threshold: 0.6
steps:

  - order: 1
    instruction: |
      ## 🎯 시나리오 1 (20점) — fw nftables drop 룰 + tcpdump 검증

      attacker IP (10.20.30.202) 의 fw input chain 진입을 drop + tcpdump 로 SYN 캡처
      확인 + cleanup.

      ## 💻 실행 명령

      ```bash
      # 1. drop 룰 추가
      ssh 6v6-fw 'sudo nft insert rule inet six_filter input position 0 ip saddr 10.20.30.202 counter log prefix "EXAM-DROP: " drop'

      # 2. tcpdump (다른 터미널)
      ssh 6v6-fw 'sudo timeout 5 tcpdump -ni eth0 host 10.20.30.202 -c 3' &

      # 3. attacker 시도 (timeout 예상)
      sleep 1
      ssh 6v6-attacker 'timeout 4 curl -s -o /dev/null -w "%{http_code} %{time_total}s\n" -H "Host: juice.6v6.lab" http://10.20.30.1/'
      wait

      # 4. counter 확인
      ssh 6v6-fw 'sudo nft list chain inet six_filter input | grep -B1 EXAM-DROP'

      # 5. Cleanup
      HANDLE=$(ssh 6v6-fw 'sudo nft -a list chain inet six_filter input | grep EXAM-DROP | grep -oE "handle [0-9]+" | head -1 | awk "{print \$2}"')
      ssh 6v6-fw "sudo nft delete rule inet six_filter input handle $HANDLE"
      ```

      ## 📋 평가 (각 5점)

      1. drop 룰 추가 syntax 정확
      2. tcpdump 가 SYN 캡처 (3건)
      3. attacker timeout (000 + 4초)
      4. counter packets 3+ + dmesg log prefix
      5. cleanup 완료 + 정상 복귀 (200)
    hint: nft insert + tcpdump + curl + cleanup
    category: exam
    points: 20
    answer: |
      ssh 6v6-fw 'sudo nft insert rule inet six_filter input position 0 ip saddr 10.20.30.202 counter log prefix "EXAM-DROP: " drop'
      ssh 6v6-attacker 'timeout 4 curl -s -o /dev/null -w "%{http_code}\n" -H "Host: juice.6v6.lab" http://10.20.30.1/'
      ssh 6v6-fw 'sudo nft list chain inet six_filter input | grep EXAM-DROP'
      HANDLE=$(ssh 6v6-fw 'sudo nft -a list chain inet six_filter input | grep EXAM-DROP | grep -oE "handle [0-9]+" | head -1 | awk "{print \$2}"')
      ssh 6v6-fw "sudo nft delete rule inet six_filter input handle $HANDLE"
    answer_detail: 시나리오 1 의 5 평가 항목 모두 수행.
    verify:
      type: output_contains
      expect: '000\|200'
      field: stdout
      semantic:
        intent: "fw nftables drop 룰 종합"
        success_criteria: ["drop 효과", "counter", "cleanup"]
        acceptable_methods:
          - "**[OSS] nft insert position 0**"
          - "**[OSS] tcpdump 캡처**"
          - "**[OSS] counter packets**"
          - "**[OSS] handle 식별 + delete**"
    target_vm: bastion

  - order: 2
    instruction: |
      ## 🎯 시나리오 2 (20점) — Suricata 새 alert 룰 + 트리거

      sqlmap UA + URL `/admin` 매칭 룰 작성 + 트리거 + alert 카운트 확인 + cleanup.

      ## 💻 실행 명령

      ```bash
      ssh 6v6-ips 'sudo bash -c "echo \"alert http any any -> any any (msg:\\\"EXAM sqlmap admin\\\"; http.user_agent; content:\\\"sqlmap\\\"; nocase; http.uri; content:\\\"/admin\\\"; nocase; sid:9008001; rev:1;)\" > /etc/suricata/rules/local.rules"'
      ssh 6v6-ips 'sudo grep -q "local.rules" /etc/suricata/suricata.yaml || sudo sed -i "/rule-files:/a\\  - /etc/suricata/rules/local.rules" /etc/suricata/suricata.yaml'
      ssh 6v6-ips 'sudo suricatasc -c reload-rules'
      sleep 2
      ssh 6v6-attacker 'curl -s -o /dev/null -A "sqlmap/1.5" -H "Host: juice.6v6.lab" http://10.20.30.1/admin'
      sleep 4
      ssh 6v6-ips 'sudo grep "9008001" /var/log/suricata/eve.json | wc -l'
      ssh 6v6-ips 'sudo rm -f /etc/suricata/rules/local.rules; sudo suricatasc -c reload-rules'
      ```

      ## 📋 평가 (각 5점)

      1. 룰 syntax (sid 9008001 + http.user_agent + http.uri)
      2. reload OK
      3. alert 1+ 발생
      4. cleanup
    hint: Suricata 룰 + reload + 트리거
    category: exam
    points: 20
    answer: |
      ssh 6v6-ips 'sudo bash -c "echo \"alert http any any -> any any (msg:\\\"EXAM sqlmap admin\\\"; http.user_agent; content:\\\"sqlmap\\\"; nocase; http.uri; content:\\\"/admin\\\"; nocase; sid:9008001; rev:1;)\" > /etc/suricata/rules/local.rules"'
      ssh 6v6-ips 'sudo grep -q "local.rules" /etc/suricata/suricata.yaml || sudo sed -i "/rule-files:/a\\  - /etc/suricata/rules/local.rules" /etc/suricata/suricata.yaml'
      ssh 6v6-ips 'sudo suricatasc -c reload-rules'
      sleep 2
      ssh 6v6-attacker 'curl -s -o /dev/null -A "sqlmap/1.5" -H "Host: juice.6v6.lab" http://10.20.30.1/admin'
      sleep 4
      ssh 6v6-ips 'sudo grep "9008001" /var/log/suricata/eve.json | wc -l'
      ssh 6v6-ips 'sudo rm -f /etc/suricata/rules/local.rules; sudo suricatasc -c reload-rules'
    answer_detail: 시나리오 2.
    verify:
      type: output_contains
      expect: '\d'
      field: stdout
      semantic:
        intent: "Suricata 룰 작성"
        success_criteria: ["reload OK", "alert"]
        acceptable_methods: ["**[OSS] http.user_agent + http.uri 결합**"]
    target_vm: ips

  - order: 3
    instruction: |
      ## 🎯 시나리오 3 (20점) — ModSec 4 공격 차단 검증

      XSS / SQLi / LFI / RCE 4 공격 모두 403 + 매치 룰 ID 추출.

      ## 💻 실행 명령

      ```bash
      for payload in "?q=<script>alert(1)</script>" "?q=1'+OR+'1'='1" "?file=../../etc/passwd" "?cmd=;ls"; do
        code=$(ssh 6v6-attacker "curl -s -o /dev/null -w '%{http_code}' -H 'Host: juice.6v6.lab' \"http://10.20.30.1/$payload\"")
        echo "$code $payload"
      done
      sleep 3
      ssh 6v6-web 'sudo tail -30 /var/log/apache2/modsec_audit.log | jq -r ".audit_data.messages[]?" 2>/dev/null | grep -oE ''id "[0-9]+"'' | sort -u | head -10'
      ```

      ## 📋 평가 (각 5점)

      1. 4 payload 모두 403
      2. 941 (XSS) 매치
      3. 942 (SQLi) 매치
      4. 930 (LFI) + 932 (RCE) 매치
    hint: 4 공격 + 룰 매트릭스
    category: exam
    points: 20
    answer: |
      for payload in '?q=<script>alert(1)</script>' "?q=1'+OR+'1'='1" "?file=../../etc/passwd" "?cmd=;ls"; do
        code=$(ssh 6v6-attacker "curl -s -o /dev/null -w '%{http_code}' -H 'Host: juice.6v6.lab' \"http://10.20.30.1/$payload\"")
        echo "$code $payload"
      done
      sleep 3
      ssh 6v6-web 'sudo tail -30 /var/log/apache2/modsec_audit.log | jq -r ".audit_data.messages[]?" 2>/dev/null | grep -oE ''id "[0-9]+"'' | sort -u | head -10'
    answer_detail: 시나리오 3.
    verify:
      type: output_contains
      expect: '403'
      field: stdout
      semantic:
        intent: "4 OWASP 차단"
        success_criteria: ["403", "941+942 룰"]
        acceptable_methods: ["**[OSS] OWASP A03 / A05 / A10 매트릭스**"]
    target_vm: bastion

  - order: 4
    instruction: |
      ## 🎯 시나리오 4 (20점) — osquery 4 헌팅 쿼리

      web 호스트에서 4 헌팅 쿼리 (on_disk=0 / world-writable / cron / listening).

      ## 💻 실행 명령

      ```bash
      ssh 6v6-web 'sudo osqueryi --json "SELECT pid, name FROM processes WHERE on_disk = 0;"'
      ssh 6v6-web 'sudo osqueryi --json "SELECT path FROM file WHERE path LIKE '"'"'/tmp/%'"'"' AND mode = '"'"'0777'"'"' LIMIT 3;"'
      ssh 6v6-web 'sudo osqueryi --json "SELECT command FROM crontab LIMIT 5;"'
      ssh 6v6-web 'sudo osqueryi --json "SELECT l.port, p.name FROM listening_ports l JOIN processes p ON l.pid=p.pid;"'
      ```

      ## 📋 평가 (각 5점)

      1. on_disk=0 (빈 결과 = 정상)
      2. world-writable (빈 결과 또는 system 임시)
      3. crontab 정상 entry
      4. listening_ports JOIN 결과
    hint: osquery 헌팅 4종
    category: exam
    points: 20
    answer: |
      ssh 6v6-web 'sudo osqueryi --json "SELECT pid, name FROM processes WHERE on_disk = 0;"'
      ssh 6v6-web 'sudo osqueryi --json "SELECT path FROM file WHERE path LIKE '"'"'/tmp/%'"'"' AND mode = '"'"'0777'"'"' LIMIT 3;"'
      ssh 6v6-web 'sudo osqueryi --json "SELECT command FROM crontab LIMIT 5;"'
      ssh 6v6-web 'sudo osqueryi --json "SELECT l.port, p.name FROM listening_ports l JOIN processes p ON l.pid=p.pid;"'
    answer_detail: 시나리오 4.
    verify:
      type: output_contains
      expect: '\['
      field: stdout
      semantic:
        intent: "osquery 헌팅"
        success_criteria: ["4 쿼리 실행"]
        acceptable_methods: ["**[OSS] 4 헌팅 표준**"]
    target_vm: web

  - order: 5
    instruction: |
      ## 🎯 시나리오 5 (20점) — 통합 R/B/P 보고서

      위 4 시나리오의 결과를 종합 + Coverage Matrix 작성.

      ## 💻 산출물

      ```markdown
      # W08 중간고사 R/B/P 통합 보고서

      | 시나리오 | 결과 | 매치 |
      | 1 fw drop | pass | counter 3 |
      | 2 Suricata 룰 | pass | sid 9008001 |
      | 3 ModSec 4 공격 | pass | 4 × 403 |
      | 4 osquery 헌팅 | pass | 4 쿼리 |

      ## 발견 / 운영 권장
      ```
    hint: 종합 보고서
    category: exam
    points: 20
    answer: 종합 보고서
    answer_detail: 시나리오 5.
    verify:
      type: output_contains
      expect: W08
      field: stdout
      semantic:
        intent: "종합 R/B/P"
        success_criteria: ["4 시나리오 매트릭스"]
        acceptable_methods: ["**[관리적]** 종합 보고서"]
    target_vm: bastion
