# contents/labs/battle-adv-ai/week08.yaml
lab_id: battle-adv-ai-week08
title: 포렌식 분석 — 메모리/디스크/네트워크 증거 수집과 타임라인 재구성 (AI 지원)
version: ai
course: battle-adv-ai
week: 8
description: 디지털 포렌식 기법을 활용한 침해 사고 조사를 실습한다. 메모리 아티팩트 수집, 디스크 타임라인 분석, 숨겨진 프로세스 탐지,
  IOC 추출, 증거 보존(Chain of Custody)을 수행한다. AI SubAgent가 자동으로 명령을 실행하고 결과를 검증합니다.
difficulty: hard
duration_minutes: 150
objectives:
- 메모리 포렌식과 라이브 시스템 분석을 수행한다
- 디스크 아티팩트 수집과 타임라인을 재구성한다
- 악성코드 정적 분석과 엔트로피 분석을 실습한다
- IOC를 추출하고 증거 보존 절차를 수행한다
prerequisites:
- battle-adv-nonai week07 이수
infra_requirements:
- attacker (10.20.30.201) — nmap, metasploit, sqlmap, hydra, nikto
- secu (10.20.30.1) — nftables, suricata
- web (10.20.30.80) — apache, modsecurity, JuiceShop(:3000), DVWA(:8080)
- siem (10.20.30.100) — wazuh-manager
pass_threshold: 0.6
steps:
- order: 1
  instruction: '[BLUE] 메모리 포렌식 아티팩트를 수집하라. /proc 파일시스템에서 모든 프로세스의 실행 경로, 명령줄, 네트워크
    연결을 수집하라.'
  hint: /proc/[pid]/exe, cmdline 파일을 읽어 프로세스 정보를 수집한다.
  category: analysis
  points: 10
  answer: 'echo ''=== Process Forensic Analysis ==='' && for pid in $(ls /proc/ |
    grep -E ''^[0-9]+$'' | head -20); do exe=$(readlink /proc/$pid/exe 2>/dev/null)
    && cmdline=$(tr ''\0'' '' '' < /proc/$pid/cmdline 2>/dev/null) && echo "PID:$pid
    EXE:$exe CMD:${cmdline:0:80}"; done | head -15 && ps auxf > /tmp/forensic_processes.txt
    && ss -tnap > /tmp/forensic_connections.txt && echo ''Artifacts collected: processes,
    connections'''
  answer_detail: 라이브 메모리 포렌식에서 /proc 파일시스템은 프로세스 정보의 보고이다. exe(실행 파일 경로), cmdline(명령줄
    인수)을 분석하여 의심 프로세스를 식별한다.
  verify:
    type: output_contains
    expect: null
  target_vm: secu
  script: 'echo ''=== Process Forensic Analysis ==='' && for pid in $(ls /proc/ |
    grep -E ''^[0-9]+$'' | head -20); do exe=$(readlink /proc/$pid/exe 2>/dev/null)
    && cmdline=$(tr ''\0'' '' '' < /proc/$pid/cmdline 2>/dev/null) && echo "PID:$pid
    EXE:$exe CMD:${cmdline:0:80}"; done | head -15 && ps auxf > /tmp/forensic_processes.txt
    && ss -tnap > /tmp/forensic_connections.txt && echo ''Artifacts collected: processes,
    connections'''
  risk_level: low
  bastion_prompt: '[BLUE] 메모리 포렌식 아티팩트를 수집하라. /proc 파일시스템에서 모든 프로세스의 실행 경로, 명령줄, 네트워크
    연결을 수집해줘'
- order: 2
  instruction: '[BLUE] 숨겨진 프로세스를 탐지하라. /proc에서 보이는 PID와 ps 출력을 비교하여 루트킷에 의해 숨겨진 프로세스가
    있는지 확인하라.'
  hint: /proc의 PID 디렉토리와 ps aux의 PID를 비교한다.
  category: analysis
  points: 10
  answer: 'echo "=== Hidden Process Detection ===" && PROC_PIDS=$(ls -d /proc/[0-9]*
    2>/dev/null | sed "s|/proc/||" | sort -n) && PS_PIDS=$(ps -eo pid --no-headers
    | tr -d " " | sort -n) && HIDDEN=$(comm -23 <(echo "$PROC_PIDS") <(echo "$PS_PIDS"))
    && if [ -z "$HIDDEN" ]; then echo "No hidden processes detected - system appears
    clean"; else echo "Hidden PIDs: $HIDDEN"; fi'
  answer_detail: 루트킷은 ps 등의 도구를 후킹하여 프로세스를 숨긴다. /proc(커널 직접 노출)과 ps(사용자 공간)의 PID를
    비교하면 불일치를 발견할 수 있다.
  verify:
    type: output_contains
    expect: null
  target_vm: secu
  script: 'echo "=== Hidden Process Detection ===" && PROC_PIDS=$(ls -d /proc/[0-9]*
    2>/dev/null | sed "s|/proc/||" | sort -n) && PS_PIDS=$(ps -eo pid --no-headers
    | tr -d " " | sort -n) && HIDDEN=$(comm -23 <(echo "$PROC_PIDS") <(echo "$PS_PIDS"))
    && if [ -z "$HIDDEN" ]; then echo "No hidden processes detected - system appears
    clean"; else echo "Hidden PIDs: $HIDDEN"; fi'
  risk_level: low
  bastion_prompt: '[BLUE] 숨겨진 프로세스를 탐지하라. /proc에서 보이는 PID와 ps 출력을 비교하여 루트킷에 의해 숨겨진
    프로세스가 있는지 확인해줘'
- order: 3
  instruction: '[BLUE] 디스크 아티팩트를 수집하라. 최근 24시간 내 수정된 파일, /tmp의 스크립트/바이너리, SUID 바이너리를
    조사하라.'
  hint: find 명령으로 최근 변경 파일을 검색한다.
  category: analysis
  points: 10
  answer: echo '=== Disk Artifact Collection ===' && echo '[1] Recently modified (24h):'
    && find / -mtime -1 -type f -not -path '/proc/*' -not -path '/sys/*' -not -path
    '/run/*' 2>/dev/null | head -15 && echo '[2] Suspicious in /tmp:' && find /tmp
    -type f -name '*.sh' -o -name '*.py' -o -name '*.elf' -o -name '*.bin' 2>/dev/null
    | head -10 && echo '[3] SUID binaries:' && find /usr -perm -4000 -type f 2>/dev/null
    | head -10
  answer_detail: 최근 수정된 파일(mtime), 임시 디렉토리의 스크립트/바이너리, SUID 비트 설정 파일은 핵심 포렌식 아티팩트이다.
  verify:
    type: output_contains
    expect: null
  target_vm: secu
  script: echo '=== Disk Artifact Collection ===' && echo '[1] Recently modified (24h):'
    && find / -mtime -1 -type f -not -path '/proc/*' -not -path '/sys/*' -not -path
    '/run/*' 2>/dev/null | head -15 && echo '[2] Suspicious in /tmp:' && find /tmp
    -type f -name '*.sh' -o -name '*.py' -o -name '*.elf' -o -name '*.bin' 2>/dev/null
    | head -10 && echo '[3] SUID binaries:' && find /usr -perm -4000 -type f 2>/dev/null
    | head -10
  risk_level: low
  bastion_prompt: '[BLUE] 디스크 아티팩트를 수집해줘'
- order: 4
  instruction: '[BLUE] 파일 시스템 타임라인을 생성하라. /tmp과 /var/log의 MAC 타임스탬프를 시간순으로 정렬하라.'
  hint: find -printf로 타임스탬프를 추출하고 정렬한다.
  category: analysis
  points: 10
  answer: echo "=== File System Timeline ===" && find /tmp /var/log -maxdepth 2 -type
    f -printf "%T+ %p\n" 2>/dev/null | sort -r | head -20 && echo "---" && echo "Focus
    on clusters of modifications within short timeframes"
  answer_detail: 파일 시스템 타임라인은 MAC 타임스탬프를 시간순으로 정렬하여 공격 활동을 재구성한다. 짧은 시간 내 대량의 파일 수정
    클러스터는 공격 행위를 나타낸다.
  verify:
    type: output_contains
    expect: null
  target_vm: secu
  script: echo "=== File System Timeline ===" && find /tmp /var/log -maxdepth 2 -type
    f -printf "%T+ %p\n" 2>/dev/null | sort -r | head -20 && echo "---" && echo "Focus
    on clusters of modifications within short timeframes"
  risk_level: low
  bastion_prompt: '[BLUE] 파일 시스템 타임라인을 생성하라. /tmp과 /var/log의 MAC 타임스탬프를 시간순으로 정렬해줘'
- order: 5
  instruction: '[BLUE] 악성코드 정적 분석을 시뮬레이션하라. 의심 파일의 해시, 문자열, ELF 헤더를 분석하고 파일 엔트로피를
    계산하라.'
  hint: sha256sum, strings, file, readelf로 바이너리를 분석한다.
  category: analysis
  points: 10
  answer: "echo '=== Static Analysis + Entropy ===' && SAMPLE=/usr/bin/curl && sha256sum\
    \ $SAMPLE && file $SAMPLE && strings $SAMPLE | grep -iE 'http|socket|connect|exec'\
    \ | head -5 && python3 -c \"\nimport math, os\nfrom collections import Counter\n\
    def entropy(data):\n    if not data: return 0\n    freq = Counter(data)\n    length\
    \ = len(data)\n    return -sum((c/length) * math.log2(c/length) for c in freq.values())\n\
    files = ['/usr/bin/ls', '/usr/bin/curl', '/etc/passwd']\nprint('File Entropy Analysis:')\n\
    for f in files:\n    try:\n        with open(f, 'rb') as fh:\n            data\
    \ = fh.read(10000)\n        e = entropy(data)\n        assess = 'PACKED/ENC' if\
    \ e > 7.0 else 'NORMAL'\n        print(f'  {f}: entropy={e:.4f} -> {assess}')\n\
    \    except: pass\n\""
  answer_detail: 정적 분석은 악성코드를 실행하지 않고 분석한다. Shannon 엔트로피 7.0 이상은 패킹/암호화된 악성코드 가능성이
    높다.
  verify:
    type: output_contains
    expect: null
  target_vm: secu
  script: "echo '=== Static Analysis + Entropy ===' && SAMPLE=/usr/bin/curl && sha256sum\
    \ $SAMPLE && file $SAMPLE && strings $SAMPLE | grep -iE 'http|socket|connect|exec'\
    \ | head -5 && python3 -c \"\nimport math, os\nfrom collections import Counter\n\
    def entropy(data):\n    if not data: return 0\n    freq = Counter(data)\n    length\
    \ = len(data)\n    return -sum((c/length) * math.log2(c/length) for c in freq.values())\n\
    files = ['/usr/bin/ls', '/usr/bin/curl', '/etc/passwd']\nprint('File Entropy Analysis:')\n\
    for f in files:\n    try:\n        with open(f, 'rb') as fh:\n            data\
    \ = fh.read(10000)\n        e = entropy(data)\n        assess = 'PACKED/ENC' if\
    \ e > 7.0 else 'NORMAL'\n        print(f'  {f}: entropy={e:.4f} -> {assess}')\n\
    \    except: pass\n\""
  risk_level: low
  bastion_prompt: '[BLUE] 악성코드 정적 분석을 시뮬레이션해줘'
- order: 6
  instruction: '[BLUE] 증거 보존(Chain of Custody)을 수행하라. 수집된 모든 포렌식 아티팩트의 SHA256 해시를
    계산하고 증거 목록을 작성하라.'
  hint: sha256sum으로 모든 증거 파일의 해시를 기록한다.
  category: response
  points: 10
  answer: 'echo "=== Chain of Custody Record ===" > /tmp/evidence_custody.txt && echo
    "Investigator: Lab Student" >> /tmp/evidence_custody.txt && echo "Date: $(date
    -u +%Y-%m-%dT%H:%M:%SZ)" >> /tmp/evidence_custody.txt && echo "Case ID: FORENSIC-$(date
    +%Y%m%d)-001" >> /tmp/evidence_custody.txt && echo "" >> /tmp/evidence_custody.txt
    && for f in /tmp/forensic_*.txt; do if [ -f "$f" ]; then HASH=$(sha256sum "$f"
    | awk "{print \$1}"); SIZE=$(stat -c%s "$f"); echo "File: $f | Size: $SIZE | SHA256:
    $HASH" >> /tmp/evidence_custody.txt; fi; done && cat /tmp/evidence_custody.txt'
  answer_detail: Chain of Custody는 법적 절차에서 증거 무결성을 보장한다. SHA256 해시, 수집 시각, 수집자, 케이스
    ID를 기록한다.
  verify:
    type: output_contains
    expect: null
  target_vm: secu
  script: 'echo "=== Chain of Custody Record ===" > /tmp/evidence_custody.txt && echo
    "Investigator: Lab Student" >> /tmp/evidence_custody.txt && echo "Date: $(date
    -u +%Y-%m-%dT%H:%M:%SZ)" >> /tmp/evidence_custody.txt && echo "Case ID: FORENSIC-$(date
    +%Y%m%d)-001" >> /tmp/evidence_custody.txt && echo "" >> /tmp/evidence_custody.txt
    && for f in /tmp/forensic_*.txt; do if [ -f "$f" ]; then HASH=$(sha256sum "$f"
    | awk "{print \$1}"); SIZE=$(stat -c%s "$f"); echo "File: $f | Size: $SIZE | SHA256:
    $HASH" >> /tmp/evidence_custody.txt; fi; done && cat /tmp/evidence_custody.txt'
  risk_level: low
  bastion_prompt: secu VM(보안게이트웨이)에서 cat으로 [BLUE] 증거 보존(Chain of Custody)을 수행 작업을
    수행해줘
- order: 7
  instruction: '[BLUE] IOC(Indicator of Compromise)를 STIX 형식으로 문서화하라.'
  hint: IP, 도메인, 해시, 파일 경로 등의 IOC를 STIX indicator로 변환한다.
  category: analysis
  points: 10
  answer: cat > /tmp/forensic_ioc.json << 'EOF' {"type":"bundle","id":"bundle--forensic-001","objects":[{"type":"indicator","id":"indicator--f001","name":"Suspicious
    crontab modification","pattern":"[file:name = 'crontabs' AND file:parent_directory_ref.path
    = '/var/spool/cron']","pattern_type":"stix","indicator_types":["malicious-activity"],"valid_from":"2026-04-06T00:00:00Z"},{"type":"indicator","id":"indicator--f002","name":"DNS
    tunneling domain","pattern":"[domain-name:value LIKE '%.exfil.test.local']","pattern_type":"stix","indicator_types":["malicious-activity"],"valid_from":"2026-04-06T00:00:00Z"}]}
    EOF python3 -m json.tool /tmp/forensic_ioc.json | head -15
  answer_detail: 포렌식에서 추출된 IOC를 STIX 형식으로 문서화하면 다른 조직/도구와 공유할 수 있다.
  verify:
    type: output_contains
    expect: null
  target_vm: siem
  script: cat > /tmp/forensic_ioc.json << 'EOF' {"type":"bundle","id":"bundle--forensic-001","objects":[{"type":"indicator","id":"indicator--f001","name":"Suspicious
    crontab modification","pattern":"[file:name = 'crontabs' AND file:parent_directory_ref.path
    = '/var/spool/cron']","pattern_type":"stix","indicator_types":["malicious-activity"],"valid_from":"2026-04-06T00:00:00Z"},{"type":"indicator","id":"indicator--f002","name":"DNS
    tunneling domain","pattern":"[domain-name:value LIKE '%.exfil.test.local']","pattern_type":"stix","indicator_types":["malicious-activity"],"valid_from":"2026-04-06T00:00:00Z"}]}
    EOF python3 -m json.tool /tmp/forensic_ioc.json | head -15
  risk_level: low
  bastion_prompt: '[BLUE] IOC(Indicator of Compromise)를 STIX 형식으로 문서화해줘'
- order: 8
  instruction: '[RED/BLUE] 환경을 정리하고 침해 사고 포렌식 보고서를 NIST SP 800-86 형식으로 작성하라.'
  hint: Evidence, Analysis Methods, Findings, Recovery Plan, Recommendations를 포함한다.
  category: analysis
  points: 15
  answer: "rm -f /tmp/forensic_*.txt /tmp/forensic_*.json /tmp/evidence_*.txt 2>/dev/null\
    \ && cat > /tmp/forensic_report.txt << 'EOF'\n=== Digital Forensic Investigation\
    \ Report ===\nCase ID: FORENSIC-20260406-001\nMethodology: NIST SP 800-86\n\n\
    1. Evidence Collected:\n   - Memory: processes, modules, network connections\n\
    \   - Disk: file system timeline, SUID binaries\n   - Log: auth.log, syslog, Suricata\
    \ eve.json\n   - Network: ARP table, routing, DNS\n\n2. Analysis:\n   - Hidden\
    \ process detection via /proc comparison\n   - File entropy for packed binary\
    \ detection\n   - Timeline integration (MAC timestamps)\n   - IOC extraction in\
    \ STIX format\n\n3. Findings:\n   [F-01] No hidden processes detected\n   [F-02]\
    \ All binaries have normal entropy\n   [F-03] Timeline shows normal file activity\n\
    \n4. Recovery Plan:\n   P1: Remove any persistence mechanisms found\n   P2: Patch\
    \ identified vulnerabilities\n   P3: Implement FIM on critical servers\n\n5. Recommendations:\n\
    \   - Deploy FIM on all servers\n   - Enable DNS anomaly detection\n   - Conduct\
    \ regular forensic readiness assessments\n\nMITRE ATT&CK: T1190, T1053, T1048,\
    \ T1005, T1070\nEOF\ncat /tmp/forensic_report.txt"
  answer_detail: NIST SP 800-86 포렌식 보고서는 증거 목록, 분석 방법, 발견 사항, 복구 계획, 권고사항을 포함한다. Chain
    of Custody를 유지하고 모든 발견에 증거를 명시한다.
  verify:
    type: output_contains
    expect: null
  target_vm: attacker
  script: "rm -f /tmp/forensic_*.txt /tmp/forensic_*.json /tmp/evidence_*.txt 2>/dev/null\
    \ && cat > /tmp/forensic_report.txt << 'EOF'\n=== Digital Forensic Investigation\
    \ Report ===\nCase ID: FORENSIC-20260406-001\nMethodology: NIST SP 800-86\n\n\
    1. Evidence Collected:\n   - Memory: processes, modules, network connections\n\
    \   - Disk: file system timeline, SUID binaries\n   - Log: auth.log, syslog, Suricata\
    \ eve.json\n   - Network: ARP table, routing, DNS\n\n2. Analysis:\n   - Hidden\
    \ process detection via /proc comparison\n   - File entropy for packed binary\
    \ detection\n   - Timeline integration (MAC timestamps)\n   - IOC extraction in\
    \ STIX format\n\n3. Findings:\n   [F-01] No hidden processes detected\n   [F-02]\
    \ All binaries have normal entropy\n   [F-03] Timeline shows normal file activity\n\
    \n4. Recovery Plan:\n   P1: Remove any persistence mechanisms found\n   P2: Patch\
    \ identified vulnerabilities\n   P3: Implement FIM on critical servers\n\n5. Recommendations:\n\
    \   - Deploy FIM on all servers\n   - Enable DNS anomaly detection\n   - Conduct\
    \ regular forensic readiness assessments\n\nMITRE ATT&CK: T1190, T1053, T1048,\
    \ T1005, T1070\nEOF\ncat /tmp/forensic_report.txt"
  risk_level: low
  bastion_prompt: attacker VM에서 cat으로 [RED/BLUE] 환경을 정리하고 침해 사고 포렌식 보고서를 NIST SP 800-86
    형식으로 작업을 수행해줘
- order: 99
  instruction: '[Multi-task] 이번 주 핵심 작업 5개를 한 번에 수행하라.'
  bastion_prompt: '다음 작업들을 순서대로 수행해줘:

    1) [BLUE] 메모리 포렌식 아티팩트를 수집하라. /proc 파일시스템에서 모든 프로세스의 실행 경로, 명령줄, 네트워크 연결을 수집해줘

    2) [BLUE] 숨겨진 프로세스를 탐지하라. /proc에서 보이는 PID와 ps 출력을 비교하여 루트킷에 의해 숨겨진 프로세스가 있는지 확인해줘

    3) [BLUE] 디스크 아티팩트를 수집해줘

    4) [BLUE] 파일 시스템 타임라인을 생성하라. /tmp과 /var/log의 MAC 타임스탬프를 시간순으로 정렬해줘

    5) [BLUE] 악성코드 정적 분석을 시뮬레이션해줘'
  hint: 복합 작업은 Bastion이 순서대로 계획하고 실행해야 한다
  category: multi_task
  points: 50
  verify:
    type: output_contains
    expect: ''
    field: stdout
  target_vm: secu
  risk_level: low
  answer: ''
  answer_detail: Bastion이 5개 작업을 순차적으로 계획(planning)하고 실행(executing)한 뒤 종합 분석(validating)하는
    다단계 작업.
  script: ''
