# contents/labs/soc-adv-nonai/week13.yaml
lab_id: soc-adv-nonai-week13
title: Purple Team 연습 - 공격과 방어의 통합
version: non-ai
course: soc-adv-nonai
week: 13
description: Red Team(공격)과 Blue Team(방어)이 협업하여 탐지 역량을 검증하고 개선하는 Purple Team 연습 실습
difficulty: hard
duration_minutes: 120
objectives:
- Purple Team 방법론을 이해하고 적용할 수 있다
- ATT&CK 기반 공격 시뮬레이션과 탐지 검증을 수행할 수 있다
- 탐지 갭을 식별하고 개선할 수 있다
pass_threshold: 0.6
steps:
- order: 1
  instruction: Purple Team의 개념과 Red/Blue Team과의 차이를 설명하라. Purple Team 연습의 목적, 참여자, 프로세스를 정의하라.
  hint: Purple Team은 대립이 아닌 협업에 초점을 맞춘다
  category: framework
  points: 8
  answer: 'Red Team: 목표 달성(침투), 독립적. Blue Team: 탐지/대응. Purple Team: 협업→공격 시뮬레이션+실시간 탐지 검증. 목적: 탐지 갭 식별과
    즉각적 개선'
  answer_detail: Purple Team은 Red Team의 공격 기법을 Blue Team의 탐지 관점에서 즉시 검증한다. 전통적 Red Team 평가와 달리 결과를 기다리지
    않고, 각 기법별로 즉시 탐지 여부를 확인하고 개선한다.
  verify:
    type: output_contains
    expect: Purple
    semantic:
      intent: Purple Team 개념 — Red Team (목표 달성 침투/독립) vs Blue Team (탐지/대응) vs Purple Team (협업 공격 시뮬레이션
        + 실시간 탐지 검증). 대립 → 협업. 즉각 개선 루프. SANS SEC599 / NIST SP 800-160 Red-Blue-Purple.
      success_criteria:
      - Purple 키워드
      - Red/Blue/Purple 차이 설명
      - 즉각 개선 루프 개념
      acceptable_methods:
      - SANS SEC599 방법론
      - Scythe Attack Framework + Purple Team
      - Cobalt Strike + Blue Team 연동
      - VECTR (Security Risk Advisors)
      negative_signs:
      - Red Team 평가와 혼동
      - 실시간 검증 개념 없음
      - 협업 강조 누락
  target_vm: siem
- order: 2
  instruction: ATT&CK 기반 공격 시나리오를 설계하라. 5개 이상의 Tactic에 걸친 공격 체인(Initial Access → Execution → Persistence
    → Lateral Movement → Exfiltration)을 설계하라.
  hint: 각 단계에 구체적인 Technique ID와 실행 방법을 명시하라
  category: planning
  points: 12
  answer: '공격 체인: T1190(웹 취약점 악용)→T1059.004(bash 실행)→T1053.003(cron 지속성)→T1021.004(SSH 횡적이동)→T1048(대체프로토콜
    유출). 각 단계의 구체적 명령 포함'
  answer_detail: 공격 시나리오는 실제 위협 행위자의 TTP를 모방해야 한다. 각 단계의 ATT&CK Technique을 명시하고, 탐지 예상 여부를 사전에 기록하여 연습
    후 비교한다.
  verify:
    type: output_contains
    expect: ATT&CK
    semantic:
      intent: ATT&CK 공격 시나리오 5 단계 체인 — T1190 (웹 취약점 악용) → T1059.004 (bash 실행) → T1053.003 (cron 지속성) →
        T1021.004 (SSH 횡적이동) → T1048 (대체 프로토콜 유출). 각 단계 Technique ID + 구체 명령 + 탐지 예상. Kill Chain 기반.
      success_criteria:
      - ATT&CK 키워드
      - 5 Tactic 이상 체인
      - Technique ID 명시
      acceptable_methods:
      - ATT&CK Navigator scenario plan
      - MITRE Caldera 자동화 + 체인
      - Atomic Red Team multi-stage
      - Scythe TTP 시뮬레이션
      negative_signs:
      - Technique ID 없음
      - 5 Tactic 미만
      - 구체 명령 없음
  target_vm: siem
- order: 3
  instruction: '공격 시뮬레이션 1단계: Initial Access를 실행하라. 웹서버({{WEB_IP}})의 JuiceShop에 SQL Injection 공격을 수행하고,
    Blue Team의 탐지 여부를 확인하라.'
  hint: OWASP JuiceShop의 알려진 SQLi 취약점을 활용하라
  category: simulation
  points: 10
  answer: '🖥️ UI: 브라우저에서 http://WEB_IP:3000 접속

    📎 참고 (CLI): 공격: JuiceShop SQLi 실행(로그인 우회 또는 데이터 추출). 탐지 확인: Wazuh 알림 발생 여부, Suricata 알림 발생 여부, WAF
    로그 기록 여부'
  answer_detail: Initial Access 단계의 탐지는 SOC의 첫 번째 방어선이다. SQLi 공격이 SIEM, IDS, WAF 중 어디에서 탐지되었는지 확인하고, 미탐지
    지점을 기록한다.
  verify:
    type: output_contains
    expect: SQL
    semantic:
      intent: Initial Access 시뮬레이션 — JuiceShop SQLi 실행 (로그인 우회 또는 데이터 추출) + Wazuh 알림 발생 여부 + Suricata
        알림 여부 + WAF 로그 기록 여부. 3개 탐지 도구 중 어디서 잡혔는지 확인. 미탐지 지점 = 개선 포인트. MITRE ATT&CK T1190.
      success_criteria:
      - SQL 키워드
      - JuiceShop 또는 공격 실행
      - 3개 탐지 도구 확인
      acceptable_methods:
      - OWASP JuiceShop 알려진 SQLi 엔드포인트
      - sqlmap --url=<target>
      - Burp Suite 수동 테스트
      - Nuclei templates/vulnerabilities/
      negative_signs:
      - 공격 없이 이론만
      - 탐지 도구 단일 확인
      - T1190 맥락 없음
  target_vm: siem
- order: 4
  instruction: '공격 시뮬레이션 2단계: Execution을 실행하라. 웹쉘 업로드 또는 명령 인젝션으로 셸 접근을 확보하고, Blue Team의 탐지 여부를 확인하라.'
  hint: 웹서버에서 반환되는 명령 실행 결과를 확인하라
  category: simulation
  points: 10
  answer: '공격: 명령 인젝션으로 id, whoami 실행. 탐지: Wazuh 프로세스 모니터링 알림, Suricata HTTP 페이로드 탐지, apache 에러 로그'
  answer_detail: Execution 단계 탐지는 침투 성공을 식별하는 핵심이다. 웹서버 프로세스(apache/nginx)에서 비정상 자식 프로세스(bash, sh, python)가
    생성되면 침해 가능성이 매우 높다.
  verify:
    type: output_contains
    expect: Execution
    semantic:
      intent: Execution 시뮬레이션 — 명령 인젝션으로 id/whoami/uname 실행 + Wazuh 프로세스 모니터링 알림 + Suricata HTTP 페이로드
        + apache 에러 로그. 웹서버 (apache/nginx) → 비정상 자식 (bash/sh/python) = 침해 지표. MITRE ATT&CK T1059 Command
        Execution.
      success_criteria:
      - Execution 키워드
      - 셸 명령 실행
      - 비정상 자식 프로세스 탐지
      acceptable_methods:
      - 명령 인젝션 payload (;id 등)
      - 웹쉘 업로드 후 실행
      - LOLBAS 기법 활용
      - Metasploit exploit/multi/http/*
      negative_signs:
      - 공격 없이 이론
      - 자식 프로세스 탐지 개념 없음
      - T1059 맥락 없음
  target_vm: siem
- order: 5
  instruction: '공격 시뮬레이션 3단계: Persistence를 설정하라. cron job 또는 SSH authorized_keys에 백도어를 설치하고, Blue Team의
    탐지 여부를 확인하라.'
  hint: crontab -e 또는 echo >> ~/.ssh/authorized_keys
  category: simulation
  points: 10
  answer: '공격: cron에 리버스쉘 등록 또는 SSH 키 추가. 탐지: Wazuh FIM(crontab 변경 탐지), SSH 설정 변경 알림, 신규 cron job 알림'
  answer_detail: Persistence 탐지는 공격자의 장기 접근을 방지한다. 파일 무결성 모니터링(FIM)이 핵심이며, crontab, authorized_keys, systemd
    서비스 파일의 변경을 감시해야 한다.
  verify:
    type: output_contains
    expect: Persistence
    semantic:
      intent: Persistence 시뮬레이션 — crontab 에 리버스쉘 등록 또는 ~/.ssh/authorized_keys 추가 + Wazuh FIM (crontab
        변경 탐지) + SSH 설정 변경 알림 + 신규 cron job 알림. FIM = Persistence 탐지 핵심. MITRE ATT&CK T1053.003 cron +
        T1098.004 SSH keys.
      success_criteria:
      - Persistence 키워드
      - cron 또는 SSH key 설치
      - FIM 변경 탐지
      acceptable_methods:
      - crontab -e + 페이로드
      - echo >> authorized_keys
      - systemd-run --unit=backdoor
      - Atomic Red Team T1053.003/T1098.004
      negative_signs:
      - Persistence 개념 없음
      - FIM 탐지 언급 없음
      - T1053 맥락 없음
  target_vm: siem
- order: 6
  instruction: '공격 시뮬레이션 4단계: Lateral Movement를 수행하라. 확보한 접근 권한으로 다른 내부 시스템({{SECU_IP}}00 등)으로 이동을 시도하라.'
  hint: SSH 키 재사용 또는 발견된 비밀번호로 다른 서버 접속 시도
  category: simulation
  points: 10
  answer: '공격: web→siem SSH 접속 시도, 내부 네트워크 스캔. 탐지: nftables 내부→내부 비정상 연결 로그, Wazuh SSH 인증 알림, Suricata
    내부 스캔 탐지'
  answer_detail: Lateral Movement는 피해 범위를 확대하는 핵심 단계이다. 내부 네트워크 모니터링이 미흡하면 탐지가 어려우므로, 내부 세그먼트 간 트래픽도 모니터링해야
    한다.
  verify:
    type: output_contains
    expect: Lateral
    semantic:
      intent: Lateral Movement 시뮬레이션 — web→siem SSH 접속 시도 + 내부 네트워크 스캔 + nftables 내부→내부 비정상 연결 로그 + Wazuh
        SSH 인증 알림 + Suricata 내부 스캔 탐지. 내부 세그먼트 간 트래픽 모니터링 중요. MITRE ATT&CK T1021 Remote Services.
      success_criteria:
      - Lateral 키워드
      - 내부 → 내부 연결 시도
      - nftables 또는 내부 트래픽 모니터링
      acceptable_methods:
      - ssh user@internal-ip + 키 재사용
      - nmap -sn 내부 서브넷
      - Impacket tools (psexec/smbexec)
      - Evil-WinRM / CrackMapExec
      negative_signs:
      - 외부 → 내부 (Initial Access와 혼동)
      - 내부 모니터링 없음
      - Lateral 개념 없음
  target_vm: siem
- order: 7
  instruction: '공격 시뮬레이션 5단계: Exfiltration을 시도하라. 수집한 데이터를 DNS 터널링, HTTP POST, 또는 기타 방법으로 외부로 전송을 시도하라.'
  hint: base64로 데이터를 인코딩하여 DNS 쿼리의 서브도메인에 포함하라
  category: simulation
  points: 10
  answer: '공격: 수집 데이터를 base64 인코딩 후 DNS 쿼리로 전송 시도. 탐지: Suricata DNS 이상 탐지, Wazuh DNS 쿼리 모니터링, 아웃바운드 트래픽
    이상'
  answer_detail: Exfiltration 탐지는 가장 어려운 영역 중 하나이다. DNS 터널링은 정상 DNS 트래픽에 숨기므로, 서브도메인 길이, 쿼리 빈도, 엔트로피 분석이
    필요하다.
  verify:
    type: output_contains
    expect: Exfiltration
    semantic:
      intent: 'Exfiltration 시뮬레이션 — base64 인코딩 후 DNS 서브도메인 포함 + Suricata DNS 이상 탐지 + Wazuh DNS 쿼리 모니터링
        + 아웃바운드 트래픽 이상. DNS 터널링 탐지 어려운 이유: 정상 DNS 트래픽에 숨김. 서브도메인 길이/빈도/엔트로피. MITRE ATT&CK T1048 + T1071.004.'
      success_criteria:
      - Exfiltration 키워드
      - DNS 터널링 또는 인코딩 유출
      - 탐지 기법 (길이/엔트로피)
      acceptable_methods:
      - iodine / dnscat2 도구
      - base64 encode + dig 쿼리
      - xxd + openssl + HTTP POST
      - Atomic Red Team T1048
      negative_signs:
      - Exfiltration 개념 없음
      - DNS 터널링 탐지 기법 없음
      - T1048 맥락 없음
  target_vm: siem
- order: 8
  instruction: 탐지 결과를 종합하라. 각 공격 단계별 탐지 성공/실패를 매트릭스로 정리하고, 탐지 도구별(Wazuh, Suricata, WAF) 기여도를 분석하라.
  hint: 공격 기법 × 탐지 도구 매트릭스를 작성하라
  category: analysis
  points: 12
  answer: '🖥️ UI: Wazuh Dashboard(https://SIEM_IP:443) 접속 → 해당 메뉴 확인

    📎 참고 (CLI): 매트릭스: Initial Access(Suricata✓, WAF✓, Wazuh○), Execution(Wazuh✓, Suricata○), Persistence(Wazuh
    FIM✓), Lateral Movement(nftables○, Wazuh✓), Exfiltration(Suricata○, Wazuh✗). ✓=탐지, ○=부분탐지, ✗=미탐지'
  answer_detail: 탐지 매트릭스는 Purple Team 연습의 핵심 산출물이다. 각 기법별로 어떤 도구가 탐지했고 어떤 도구가 실패했는지 명확히 기록하면, 즉각적인 개선
    우선순위를 도출할 수 있다.
  verify:
    type: output_contains
    expect: 매트릭스
    semantic:
      intent: '탐지 매트릭스 — 공격 기법 × 탐지 도구: Initial Access (Suricata ✓/WAF ✓/Wazuh ○) + Execution (Wazuh ✓/Suricata
        ○) + Persistence (Wazuh FIM ✓) + Lateral Movement (nftables ○/Wazuh ✓) + Exfiltration (Suricata
        ○/Wazuh ✗). ✓ 탐지 / ○ 부분 / ✗ 미탐지. 갭 시각화.'
      success_criteria:
      - 매트릭스 키워드
      - 5 단계 × 3 도구 이상
      - 탐지/부분/미탐지 3 상태
      acceptable_methods:
      - markdown 표 (Tactic | Tool | Status)
      - VECTR Purple Team tracking
      - ATT&CK Navigator layer
      - Scythe Community Threats
      negative_signs:
      - 2차원 매트릭스 없음
      - 상태 2개만 (탐지/미탐지)
      - 도구 구분 없음
  target_vm: siem
- order: 9
  instruction: 탐지 갭을 분석하라. 미탐지된 공격 기법에 대해 근본 원인(데이터 소스 부재, 룰 부재, 룰 우회)을 분석하고 개선 방안을 제시하라.
  hint: 각 미탐지에 대해 '왜 탐지하지 못했는가'를 5-Why 분석하라
  category: analysis
  points: 12
  answer: '갭 분석: 1) DNS 유출 미탐지 → DNS 쿼리 상세 로깅 미설정(데이터소스 부재), 2) 횡적이동 부분탐지 → 내부 트래픽 모니터링 룰 미흡(룰 부재). 개선:
    DNS 로그 수집 추가, 내부 횡적이동 탐지 룰 작성'
  answer_detail: '탐지 갭의 근본 원인은 크게 3가지이다: 1) 데이터가 수집되지 않음, 2) 데이터는 있지만 탐지 룰이 없음, 3) 룰은 있지만 우회됨. 각 원인에 따라
    대응 방안이 다르다.'
  verify:
    type: output_contains
    expect: 갭
    semantic:
      intent: '탐지 갭 분석 5-Why — 미탐지 기법별 근본 원인 3 분류: (1) 데이터 소스 부재 (DNS 상세 로깅 미설정) + (2) 룰 부재 (내부 횡적이동 룰
        미흡) + (3) 룰 우회 (기법이 기존 룰 회피). 각 원인별 다른 대응. Toyota 5 Whys root cause.'
      success_criteria:
      - 갭 키워드
      - 근본 원인 3 유형 분류
      - 개선 방안 매핑
      acceptable_methods:
      - 5-Why 분석법
      - Fishbone diagram (Ishikawa)
      - Root Cause Analysis (RCA) Framework
      - Blameless Postmortem
      negative_signs:
      - 갭 원인 분석 없음
      - 1 유형만 (룰 부재)
      - 대응 방안 누락
  target_vm: siem
- order: 10
  instruction: 즉시 개선 조치를 실행하라. 식별된 탐지 갭에 대해 새로운 탐지 룰을 작성하고, 동일 공격을 재실행하여 탐지를 검증하라.
  hint: 가장 위험한 갭부터 개선하라
  category: improvement
  points: 10
  answer: '개선: 1) DNS 터널링 탐지 룰 추가(서브도메인 길이 >50), 2) 내부 SSH 비정상 연결 룰 추가. 재검증: 동일 공격 재실행 → 탐지 성공 확인'
  answer_detail: 즉시 개선과 재검증은 Purple Team의 핵심 가치이다. 발견 즉시 개선하고 검증하는 루프를 통해, 연습 완료 시 모든 탐지 갭이 해소된 상태를 달성한다.
  verify:
    type: output_contains
    expect: 개선
    semantic:
      intent: '즉시 개선 + 재검증 — (1) DNS 터널링 탐지 룰 추가 (서브도메인 길이 >50) + (2) 내부 SSH 비정상 연결 룰 + 동일 공격 재실행 → 탐지
        성공 확인. Purple Team 핵심 가치: 발견 즉시 개선하고 검증하는 루프. 연습 종료 시 모든 갭 해소.'
      success_criteria:
      - 개선 키워드
      - 룰 추가 + 재검증
      - Purple Team 루프 개념
      acceptable_methods:
      - local_rules.xml 즉시 추가 + wazuh-logtest
      - Suricata custom rule + reload
      - VECTR 재테스트 기록
      - Atomic Red Team 재실행
      negative_signs:
      - 재검증 없음
      - 개선 계획만 + 즉시 실행 안 함
      - 루프 개념 없음
  target_vm: siem
- order: 11
  instruction: Atomic Red Team 프레임워크를 활용하라. 자동화된 공격 시뮬레이션 도구(Atomic Red Team)의 테스트 케이스를 실행하고 탐지를 검증하라.
  hint: 각 Atomic Test는 개별 ATT&CK Technique에 매핑된다
  category: automation
  points: 8
  answer: 'Atomic Red Team: T1053.003(cron) 테스트 실행 → 탐지 여부 확인, T1059.004(bash) 테스트 실행 → 탐지 여부 확인. 자동화된
    공격-탐지 검증 루프'
  answer_detail: Atomic Red Team은 개별 기법을 테스트하는 자동화 도구이다. 각 테스트는 독립적이므로 킬체인 전체가 아닌 개별 기법의 탐지 여부를 빠르게 검증할
    수 있다.
  verify:
    type: output_contains
    expect: Atomic
    semantic:
      intent: Atomic Red Team 자동화 — T1053.003 (cron) 테스트 + T1059.004 (bash) 테스트 개별 실행 + 탐지 여부 확인. 각 Atomic
        Test = 개별 ATT&CK Technique 매핑 + 독립 실행. 빠른 회귀 테스트. Red Canary Atomic Red Team 프레임워크.
      success_criteria:
      - Atomic 키워드
      - T 개수 2개 이상
      - 개별 기법 독립 테스트
      acceptable_methods:
      - Invoke-AtomicTest T1053.003
      - atomic-operator Python CLI
      - MITRE Caldera + atomic 통합
      - Red Canary Atomic Red Team GitHub
      negative_signs:
      - Atomic 프레임워크 사용 없음
      - 단일 테스트만
      - 자동화 맥락 없음
  target_vm: siem
- order: 12
  instruction: ATT&CK Navigator에 결과를 시각화하라. Purple Team 연습 결과를 ATT&CK Navigator layer로 생성하고, 탐지/미탐지를 색상으로
    구분하라.
  hint: 탐지=초록, 부분탐지=노랑, 미탐지=빨강, 미테스트=회색
  category: visualization
  points: 8
  answer: 'Navigator layer 생성: 테스트된 Technique별 탐지 결과를 색상 코딩. 전체 테스트 N개, 탐지 M개(초록), 부분 K개(노랑), 미탐지 L개(빨강)'
  answer_detail: ATT&CK Navigator 시각화는 경영진에게도 직관적으로 탐지 역량을 보여준다. 빨간 영역이 즉시 개선이 필요한 위험 영역이며, 시간에 따른 커버리지
    변화를 추적할 수 있다.
  verify:
    type: output_contains
    expect: Navigator
    semantic:
      intent: ATT&CK Navigator 시각화 — layer JSON 생성 + Technique 별 탐지 결과 색상 코딩 (초록 탐지 / 노랑 부분 / 빨강 미탐지 /
        회색 미테스트). 전체 N 테스트 + 탐지 M + 부분 K + 미탐지 L. 경영진에게 직관적. Red 영역 = 즉시 개선 위험 영역.
      success_criteria:
      - Navigator 키워드
      - layer JSON 생성
      - 색상 코딩 4 상태
      acceptable_methods:
      - ATT&CK Navigator https://mitre-attack.github.io/
      - DeTT&CT detect.yaml + visualization
      - VECTR ATT&CK mapping export
      - Python mitreattack-python library
      negative_signs:
      - 시각화 없이 텍스트만
      - 색상 코딩 없음
      - layer JSON 개념 없음
  target_vm: siem
- order: 13
  instruction: Purple Team 연습 프로그램을 수립하라. 분기별 연습 일정, 시나리오 선정 기준, 참여자 역할, 성과 측정 방법을 정의하라.
  hint: 위협 인텔리전스 기반으로 가장 관련 있는 위협 시나리오를 선정하라
  category: planning
  points: 8
  answer: '프로그램: 분기별 1회, 시나리오 선정(TI 기반 최신 위협), 참여자(Red:공격시뮬레이션, Blue:탐지대응, Purple Lead:조율/기록), 성과(탐지율
    향상, 갭 해소 수, MTTD 개선)'
  answer_detail: 지속적인 Purple Team 프로그램은 탐지 역량의 지속적 향상을 보장한다. 각 연습의 결과가 다음 연습의 시나리오에 반영되는 선순환 구조가 핵심이다.
  verify:
    type: output_contains
    expect: 프로그램
    semantic:
      intent: Purple Team 프로그램 수립 — 분기별 1회 + 시나리오 선정 (TI 기반 최신 위협) + 참여자 역할 (Red 공격 시뮬레이션 + Blue 탐지 대응
        + Purple Lead 조율/기록) + 성과 측정 (탐지율/갭 해소 수/MTTD 개선). 선순환 구조 (결과 → 다음 시나리오).
      success_criteria:
      - 프로그램 키워드
      - 주기 (분기별) + 참여자 역할
      - 성과 측정 방법
      acceptable_methods:
      - VECTR Purple Team lifecycle
      - SANS SEC599 Purple Team 방법론
      - MITRE ATT&CK Data Sources 기반 시나리오
      - Scythe Community Threats 카탈로그
      negative_signs:
      - 주기 없음
      - 역할 구분 없음
      - TI 기반 시나리오 선정 없음
  target_vm: siem
- order: 14
  instruction: Purple Team 연습의 성과 지표를 측정하라. 연습 전후의 탐지 커버리지, MTTD, 신규 탐지 룰 수, 개선된 로그 소스 수를 정량화하라.
  hint: 개선율(%)로 표현하면 성과가 명확해진다
  category: measurement
  points: 8
  answer: '성과: 탐지 커버리지 60%→85%(+25%p), 평균 MTTD 30분→10분(-67%), 신규 탐지 룰 15개, 추가 로그 소스 3개, 해소된 갭 12개'
  answer_detail: 정량적 성과 측정은 Purple Team 투자의 가치를 입증한다. 특히 탐지 커버리지와 MTTD 개선은 경영진에게 가장 설득력 있는 지표이다.
  verify:
    type: output_contains
    expect: 성과
    semantic:
      intent: Purple Team 성과 지표 — 탐지 커버리지 60% → 85% (+25%p) + MTTD 30분 → 10분 (-67%) + 신규 탐지 룰 15개 + 추가
        로그 소스 3개 + 해소된 갭 12개. 개선율 % 로 표현 + 경영진 설득력. 투자 가치 정량화.
      success_criteria:
      - 성과 키워드
      - 커버리지 + MTTD 수치
      - 개선율 % 표현
      acceptable_methods:
      - Before/After 비교 표
      - Grafana Purple Team dashboard
      - Excel 성과 추적
      - VECTR Metrics export
      negative_signs:
      - 수치 없이 서술
      - 개선율 없음
      - MTTD 맥락 없음
  target_vm: siem
- order: 15
  instruction: Purple Team 종합 보고서를 작성하라. 연습 시나리오, 공격/탐지 결과 매트릭스, 갭 분석, 개선 조치, 성과 지표, 향후 계획을 포함하라.
  hint: 경영진 요약에는 탐지 역량 향상의 비즈니스 가치를 포함하라
  category: reporting
  points: 14
  answer: '종합: 5단계 공격 시뮬레이션 완료, 탐지율 85%, 갭 12개 중 10개 즉시 해소, MTTD 67% 개선. 잔여 갭: DNS 유출 고도화 탐지, 클라우드 환경
    가시성. 다음 연습: APT 시나리오'
  answer_detail: Purple Team 보고서는 공격자 관점의 위험과 방어자 관점의 개선을 모두 담는다. 잔여 갭은 다음 연습 또는 도구 도입으로 해결할 로드맵과 함께 제시한다.
  verify:
    type: output_contains
    expect: Purple
    semantic:
      intent: 'Purple Team 종합 보고서 — 5단계 공격 시뮬레이션 완료 + 탐지율 85% + 갭 12개 중 10개 즉시 해소 + MTTD 67% 개선. 잔여 갭
        (DNS 유출 고도화 탐지 + 클라우드 환경 가시성) + 다음 연습 (APT 시나리오). 경영진 요약: 탐지 역량 비즈니스 가치.'
      success_criteria:
      - Purple 키워드
      - 정량 수치 + 잔여 갭
      - 다음 연습 계획
      acceptable_methods:
      - SANS Purple Team 보고 템플릿
      - Executive Summary + 기술 부록
      - ATT&CK layer 첨부
      - VECTR export + 커스텀
      negative_signs:
      - 수치 없음
      - 잔여 갭 식별 안 함
      - 다음 연습 계획 누락
  target_vm: siem
